# admins/salt-media/media/roots/dom0-music/50-add-route-to-veth.rules
KERNEL=="veth*", ACTION=="add", RUN+="/usr/bin/add-veth-route"
