# cloud/salt-formula/cloudgate/tmpl/service/yavpp-autorecovery.timer.j2
[Unit]
Description=Timer for yavpp-autorecovery

[Timer]
OnCalendar=*:{{ offset }}/{{ period }}

[Install]
WantedBy=timers.target
