# cloud/salt-formula/cloudgate/tmpl/service/yavpp-autorecovery.service.j2
[Unit]
Description='VPP Cloud Gate autorecovery'

[Service]
Type=oneshot
ExecStart=/usr/bin/env python3 -m yavpp.autorecovery --apply
