# cloud/mdb/mdb-config-salt-baremetal/initial-salt-deploy.cron
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
@reboot root /etc/cron.yandex/initial-salt-deploy.sh >/var/log/mdb-initial-salt-deploy.log 2>&1
