# admins/salt-media/kp/roots/common/files/ubic/php.init
#!/usr/bin/perl

### BEGIN INIT INFO
# Provides:          php-fpm
# Required-Start:    ubic-watchdog
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start php-fpm
# Description:       Starts PHP FCGI Process Manager with ubic
### END INIT INFO

use Ubic::Run;
