# admins/git-repo/packages/config-corba-lpythonbp-ganglia/cron.d/logster
#* * * * *	root	/usr/sbin/logster --output=ganglia AllDomainLogster /var/log/nginx/access.log 1>/dev/null 2>&1
* * * * *	root	/usr/sbin/logster --parser-option="/etc/logster/RPSTimingsLogster.conf" --output=ganglia RPSTimingsLogster /var/log/nginx/access.log 1>/dev/null 2>&1

* * * * *       root    /usr/sbin/logster --parser-option="/etc/logster/RPSTimingsLogster-afisha.conf" --output=ganglia RPSTimingsLogster /var/log/nginx/access.log 1>/dev/null 2>&1

