# cloud/mdb/salt/salt/components/postgres/conf/cron.d/check_backups
MAILTO=''
0 12 * * * root /etc/cron.yandex/check_backups.sh

