# cloud/mdb/salt/salt/components/postgres/conf/cron.d/odyssey_cores_cleanup
MAILTO=''
SHELL=/bin/bash

0 5 * * * root sleep $((RANDOM \% 3600)); cd /var/cores; for filename in ./*; do gdb odyssey "$(echo ${filename: 2})" --eval-command 't a a bt' --batch 2>&1 | grep -e 'od_system_shutdown' -e 'parser_pgoptions.c' > /dev/null 2>&1 && rm -f "$(echo ${filename: 2})" || true; done;
