# cloud/mdb/salt/salt/components/postgres/conf/cron.d/pg_archive_wals
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

*/1 * * * *   postgres flock -n -w 10 /tmp/pg_archive_wals.lock /usr/local/yandex/pg_archive_wals.sh >>/var/log/postgresql/pg_archive_wals.log 2>&1
