# cloud/mdb/salt/salt/components/postgres/conf/cron.d/pg_partman
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

30 4 * * * postgres /etc/cron.yandex/pg_partman_maintenance.py >/dev/null 2>&1
