# cloud/mdb/salt/salt/components/postgres/conf/cron.d/pg_resetup
MAILTO=''
{% from "components/postgres/pg.jinja" import pg with context %}
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
PGPASSFILE=/var/lib/postgresql/.pgpass

*/5 * * * * root sleep $((RANDOM \% 120)); flock -o -n /tmp/pg_resetup.lock timeout 86400 /etc/cron.yandex/pg_resetup.py -c {{ salt['pillar.get']('data:config:checkpoint_timeout', '5min') }} -w {{ salt['pillar.get']('data:config:wal_replay_timeout', '13h') }} -m {{ pg.version.major }} -p {{ pg.connection_pooler }} >>/var/log/postgresql/resetup.log 2>&1 || /bin/true
