# cloud/mdb/salt/salt/components/postgres/conf/cron.d/pg_recreate_unused_slot
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

*/5 * * * * postgres sleep $((RANDOM \% 120)); flock -n -w 10 /tmp/pg_recreate_unused_slot.lock /etc/cron.yandex/pg_recreate_unused_slot.py >/dev/null 2>&1
