# cloud/mdb/salt/salt/components/postgres/conf/cron.d/pg_rm_old_archive_status
{% from "components/postgres/pg.jinja" import pg with context %}
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

30 * * * *  postgres sleep $((RANDOM \% 120)); timeout -s SIGKILL 10s flock -n /tmp/rm_not_archived_xlogs.lock find {{ pg.data }}/pg_xlog/archive_status/ -name '*.ready' -mtime +2 -delete >/dev/null 2>&1
