# cloud/mdb/salt/salt/components/postgres/conf/cron.d/pg_statements_stats
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

*/5 * * * * postgres sleep $((RANDOM \% 120)); /etc/cron.yandex/clear_statements_stats.py >/dev/null 2>&1
