# cloud/mdb/salt/salt/components/postgres/conf/cron.d/wd-pgsync
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

*/2 * * * * root sleep $((RANDOM \% 60)); /etc/cron.yandex/wd_pgsync.py >/dev/null 2>&1
