# admins/cult-admin/packages/general/yandex-cult-dbs-mongoupdate/src/etc/cron.d/mongo_update
0	5	*	*	*	root	lockf -t 1 /var/lock/mongoupdate.lock /usr/bin/mongo_update.sh >/dev/null 2>&1
