# admins/salt-media/storage/ape/roots/mongo/etc/cron.d/mongobackup
13 3,14 * * *   root /bin/bash /usr/local/bin/mongo-dump.sh 2>&1

