# admins/salt-media/storage/common/roots/files/solomon/solomon.service
[Unit]
Description=Solomon agent
AssertFileIsExecutable=/usr/local/bin/solomon-agent
AssertFileNotEmpty=/etc/solomon.conf
After=multi-user.target


[Service]
User=nobody
Group=nogroup

ExecStart=/usr/local/bin/solomon-agent --config /etc/solomon.conf

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=5
SendSIGHUP=no
Restart=on-failure

[Install]
WantedBy=multi-user.target

