# cloud/ai/misc/model_server/docker/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/nginx.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760
