# cloud/ai/misc/model_server/docker/supervisor/rest-server.conf
[program:rest-server-1]
command=python3 rest_server.py -p 8000
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760

[program:rest-server-2]
command=python3 rest_server.py -p 8001
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760

[program:rest-server-3]
command=python3 rest_server.py -p 8002
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760

[program:rest-server-4]
command=python3 rest_server.py -p 8003
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760

[program:rest-server-5]
command=python3 rest_server.py -p 8004
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760

[program:rest-server-6]
command=python3 rest_server.py -p 8005
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760

[program:rest-server-7]
command=python3 rest_server.py -p 8006
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760

[program:rest-server-8]
command=python3 rest_server.py -p 8007
autostart=true
autorestart=true
stderr_logfile=/var/log/rest-server.err.log
stderr_logfile_backups=10
stderr_logfile_maxbytes=10485760
stdout_logfile=/var/log/rest-server.out.log
stdout_logfile_backups=10
stdout_logfile_maxbytes=10485760