# admins/salt-media/storage/elliptics/roots/units/elliptics-storage-files/files/etc/cron.daily/config-mds-storage
#!/bin/sh
/usr/bin/timeout -s KILL 60 /usr/bin/find /cache/defrag/ -mindepth 1 -mtime +3 -delete 2>/dev/null || /bin/true

