# cloud/mdb/salt/salt/components/jaeger-agent/conf/jaeger-agent.upstart
description "Jaeger Agent"
author "sidh"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env GOMAXPROCS=2
chdir /etc/jaeger
exec /usr/local/bin/jaeger-agent --config-file /etc/jaeger/agent.yaml

post-stop exec sleep 15
