# cloud/mdb/deploy/saltkeys/etc/upstart.conf
description "MDB deploy saltkeys"
author "sidh"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env GOMAXPROCS=2
chdir /etc/yandex/mdb-deploy-saltkeys
exec /opt/yandex/mdb-deploy-saltkeys/bin/mdb-deploy-saltkeys

post-stop exec sleep 15
