# admins/yandex-media-check-route/etc/cron.d/default-route-fixer
SHELL=/bin/bash
*/5 * * * *   root    /usr/bin/flock -xn /var/lock/default-route-fixer -c /usr/sbin/default-route-fixer &>/dev/null || /bin/true
