# cloud/blockstore/packages/yc-nbs-breakpad-systemd/debian/postinst
#!/bin/sh

set -e
case "$1" in
    configure)
        QUEUE_DIR=/var/tmp/breakpad
        mkdir -p "${QUEUE_DIR}" || :
        chown nbs "${QUEUE_DIR}" || :
        chmod 755 "${QUEUE_DIR}"  || :

        if [ -x /bin/systemctl ]; then
            systemctl daemon-reload || :
            systemctl enable nbs-breakpad-sender.service || :
            systemctl restart nbs-breakpad-sender.service || :
        else
            initctl reload-configuration || :
            stop nbs-breakpad-sender || :
            start nbs-breakpad-sender || :
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
