# cloud/mdb/salt/salt/components/mdb-search-producer/conf/supervisor.conf
[program:mdb-search-producer]
command=/opt/yandex/mdb-search-producer/bin/mdb-search-producer
environment=GOMAXPROCS=2
directory=/opt/yandex/mdb-search-producer
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=mdb-search-producer
stdout_logfile=/var/log/supervisor/mdb-search-producer-stdout.log
redirect_stderr=true
startretries=1800

