# admins/salt-media/storage/elliptics/roots/files/dom0/rc.local
#!/bin/bash
echo 16384 > /sys/class/net/br0/bridge/group_fwd_mask || exit 0

