# admins/salt-media/noc/roots/units/k8s/files/etc/cron.d/etcd-backup
16 5 * * *	root	/usr/sbin/etcd-backup >/dev/null 2>&1
