# admins/yandex-media-mysqld-watchdog/src/mysqld_watchdog.cron
* * * * * root flock -w 3 /var/run/mysqld_watchdog.lock /usr/sbin/mysqld_watchdog.sh fix
