# cloud/mdb/salt/salt/components/clickhouse/clickhouse-cleaner/clickhouse-cleaner.cron
MAILTO=''
15 0 * * * root sleep $((RANDOM \% 30)); flock -n -o /tmp/clickhouse-cleaner.lock /usr/local/yandex/clickhouse-cleaner.py 1>/dev/null 2>&1
