# cloud/mdb/salt/salt/components/mdb_cms/cron/autoduty.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

* * * * * mdb-cms flock -w 0 /tmp/mdb-cms-autoduty.lock /opt/yandex/mdb-cms/bin/mdb-cms-autoduty --config-path=/opt/yandex/mdb-cms/etc >> /var/log/mdb-cms/autoduty.log 2>&1
