# cloud/mdb/salt/salt/components/monrun2/config/monrun-cache-chown.cron
MAILTO=''
*/5 * * * * root sleep $((RANDOM \% 240)); chown -R monitor:monitor /var/cache/monrun >/dev/null 2>&1
