# cloud/mdb/salt/salt/components/monrun2/config/wd-juggler-client.cron
MAILTO=''
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * root sleep $((RANDOM \% 55)); /etc/cron.yandex/wd-juggler-client.sh >/dev/null 2>&1
