# admins/media-graphite-sender/contrib/etc/systemd/system/media-graphite-sender.service
[Unit]
Description=media-graphite-sender: forward metrics to graphite
After=network.target

[Service]
ExecStart=/usr/bin/media-graphite-sender
Restart=on-failure

[Install]
WantedBy=multi-user.target
