# admins/salt-media/kp/roots/ott/ftpserver/files/etc/xinetd.d/vsftpd
# ftp daemon.
# instances = 25 - max connection
# per_source = 5 - max connection from 1 ip
service ftp
{
     disable = no
     socket_type = stream
     wait = no
     user = root
     server = /usr/sbin/vsftpd
     flags = IPV6
     nice = 10
     per_source = 100
     instances = 2500
}

