# cloud/blockstore/packages/yc-nbs-disk-agent-systemd/init/blockstore-disk-agent.service
[Unit]
Description=Blockstore Disk Agent
After=network-online.target rc-local.service
Wants=network-online.target
ConditionPathExists=/Berkanavt/nbs-server/cfg/nbs_disk_agent.cfg
StartLimitInterval=10
StartLimitBurst=15

[Service]
Restart=always
RestartSec=1
PermissionsStartOnly=true
User=nbs
CapabilityBoundingSet=CAP_SETFCAP CAP_SYS_PTRACE CAP_SYS_ADMIN
ExecStartPre=/sbin/setcap CAP_SYS_PTRACE,CAP_SYS_ADMIN=ep /usr/bin/blockstore-disk-agent
PermissionsStartOnly=true
ExecStart=/usr/bin/yc-nbs-disk-agent-init
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=NBS_DISK_AGENT
SyslogFacility=daemon
SyslogLevel=err
LimitNOFILE=65536
LimitCORE=0
LimitMEMLOCK=32212254720

[Install]
WantedBy=multi-user.target
