# admins/salt-media/noc/roots/files/nocdev-rs/etc/rc.local
#!/bin/bash
sleep 60
/etc/network/ya-slb-tun.sh start eth0
