# admins/salt-media/storage/ape_prod/roots/front-12-prestable/etc/cron.d/cocaine-http-proxy-autofix
# vim: ft=crontab

*/3 * * * *   root zk-flock -n 6 cocs_front_autofix -c /etc/distributed-flock-media.json "bash /usr/local/bin/socketleak-autofix.sh" >/dev/null 2>&1
*/10 * * * *   root /usr/local/bin/tornado-proxy-autofix.sh

