# admins/salt-media/storage/ape_prod/roots/front-12-prestable/etc/cron.d/echo-reqtime
# vim: ft=crontab

* * * * *   root  flock -n /var/lock/check-cocaine-reqtime /bin/bash /usr/local/bin/check-cocaine-reqtime.sh >/dev/null 2>&1

