# cloud/mdb/salt/salt/components/mdb-dns/conf/supervisor.conf
[program:mdb-dns]
command=/opt/yandex/mdb-dns/bin/mdb-dns --scheme=http --host=::1 --port=26011 --config-path=/opt/yandex/mdb-dns/etc
environment=GOMAXPROCS=2
directory=/opt/yandex/mdb-dns
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=mdb-dns
stdout_logfile=/var/log/supervisor/mdb-dns.log
redirect_stderr=true
