# cloud/mdb/images/dataproc14/bootstrap/salt/components/agent/conf/supervisor.conf
[program:dataproc-agent]
command=/opt/yandex/dataproc-agent/bin/dataproc-agent --config-path=/opt/yandex/dataproc-agent/etc
environment=GOMAXPROCS=2
directory=/opt/yandex/dataproc-agent
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
stdout_logfile=/var/log/supervisor/dataproc-agent.log
redirect_stderr=true
