# admins/git-repo/packages/database-kinopoisk-sphinx/debian/database-kinopoisk-sphinx.cron.d
15 3 * * * root /usr/bin/indexer --config /etc/sphinxsearch/sphinx.conf --rotate --quiet --all
