# admins/git-repo/packages/database-kinopoisk-sphinx/debian/postinst
#!/bin/bash

update-rc.d sphinxsearch defaults
update-rc.d sphinxsearch enable
