# admins/salt-media/content/roots/mongodb_cgroups/files/mongodb_init.conf
# Ubuntu upstart file at /etc/init/mongodb.conf
{% set rsname=salt['cmd.run']('if [ -e /etc/mongodb.conf ]; then grep replSet /etc/mongodb.conf | cut -d "=" -f 2 | tr -d " " 2>/dev/null; else echo "nors"; fi ') %}
pre-start script
    for dir in /var/lib/mongodb/ /var/log/mongodb/ /var/run/mongodb/ /opt/mongodb; do 
        [ ! -d $dir ] && mkdir -p $dir
        chown -R mongodb:mongodb $dir
    done
    cgcreate -f 750 -d 750 -a root:root -t mongodb:mongodb -g cpu,memory:mongodb/{{rsname}}
    cgset -r cpu.shares=50 mongodb/{{rsname}}
    cgset -r memory.soft_limit_in_bytes={{ ((grains.mem_total|int) * 0.7)|round|int }}M mongodb/{{rsname}}
    cgset -r memory.limit_in_bytes={{ ((grains.mem_total|int) * 0.8)|round|int }}M mongodb/{{rsname}}

end script

post-start script
    sleep 1
    cat /var/run/mongodb/mongodb.pid > /sys/fs/cgroup/cpu/mongodb/{{rsname}}/tasks
    cat /var/run/mongodb/mongodb.pid > /sys/fs/cgroup/memory/mongodb/{{rsname}}/tasks
    echo 1 > /sys/fs/cgroup/memory/mongodb/{{rsname}}/memory.oom_control 
end script

start on started cgconfig
stop on runlevel [06]

script
  ENABLE_MONGODB="yes"
  #Default options
  if [ -f /etc/default/mongodb ]; then . /etc/default/mongodb; fi

  #Check defragmentation flag (from dbmanager):
  if [ -f /var/lib/dbmanager/mongo-refill-from-hidden ]; then
    echo "Mongo refill from hidden is in progress."
    exit 1;
  fi

  #Check foreign mongodb.conf
  CONFIG=/etc/mongodb-default.conf
  if [ -e /etc/mongodb.conf ]; then
    CONFIG=/etc/mongodb.conf
  fi;
  #Don't start default mongodb if we have mongos.config. Let's start mongos!
  if [ "$CONFIG" == "/etc/mongodb-default.conf" ]; then
  if [ -e /etc/mongos.conf ]; then 
    if ! `cat /etc/mongodb-default.conf | egrep ^port>/dev/null`; then
        if ! `cat /etc/mongos.conf | grep port>/dev/null`; then
            echo "Mongos config with default port detected";
            exit 0; 
        fi;
     fi;
  fi;
  fi;
  #Mongodb sometimes losts his pid
  dbpath=`cat $CONFIG | egrep ^dbpath | awk -F '=' '{print $2}'`;
  if [ -e $dbpath/mongod.lock ] ; then
        pid=`cat $dbpath/mongod.lock`; 
        cmd=`cat /proc/$pid/cmdline | tr "\000" "\n"|head -n 1 |cut -d : -f 1`
        if [ "$cmd" == "/usr/bin/mongod" ]; then 
            echo "Mongodb already start"; 
            exit 1; 
        else 
            echo "Hmm..old pid.. i remove it"; 
            rm  $dbpath/mongod.lock; 
        fi; 
  fi

  if [ "x$ENABLE_MONGODB" = "xyes" ];then
      exec  start-stop-daemon --start --pidfile /var/run/mondodb.pid \
      --chuid mongodb --exec /usr/bin/numactl -- \
          --interleave=all /usr/bin/mongod $SERVER_PARAMS --config $CONFIG
  fi
end script

post-stop script
        cgdelete cpu,memory:mongodb/{{rsname}}
end script
