# cloud/mdb/salt/salt/components/pushclient/conf/logrotate.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

*/10 * * * * root sleep $((RANDOM \% 120)); /usr/sbin/logrotate -s /var/lib/logrotate/pushclient-status /etc/pushclient/logrotate.conf
