# cloud/mdb/salt/salt/components/pushclient/conf/wd-pushclient.cron
MAILTO=''
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

* * * * * root sleep $((RANDOM \% 45)); flock -n /tmp/wd-pushclient.lock timeout 30 /etc/cron.yandex/wd-pushclient.py >/dev/null 2>&1
