# cloud/mdb/salt/salt/components/greenplum/conf/cron.d/disk_usage_watcher
MAILTO=''
*/1 * * * * gpadmin sleep $((RANDOM \% 20)); /usr/local/yandex/gp_is_master && flock -n /tmp/disk_usage_watcher.lock /usr/bin/timeout 30 /usr/local/yandex/disk_usage_watcher.py --soft {{ salt['pillar.get']('data:greenplum:disk_usage_watcher_limit', 95) }} -p /var/lib/greenplum 2>&1 | logger -t 'disk_usage_watcher'
