# cloud/mdb/salt/salt/components/greenplum/conf/cron.d/orphans_killer
{% from "components/greenplum/map.jinja" import gpdbvars with context %}

PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

# m h dom mon dow user command
*/5 * * * * {{ gpdbvars.gpadmin }} /usr/local/yandex/gp_is_master && flock -n /tmp/orphans_killer.lock /usr/local/yandex/orphans_killer.py --dry-run >> {{ gpdbvars.gplog }}/orphans_killer.log 2>&1
