# cloud/mdb/salt/salt/components/greenplum/conf/cron.d/wd-pxf
MAILTO=''
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

* * * * * root flock -n /tmp/wd-pxf.lock timeout 60 /etc/cron.yandex/wd-pxf.sh >/dev/null 2>&1
