# admins/git-repo/packages/yandex-graphite-checks-system/debian/cron.d
* * * * * root  /usr/lib/yandex-graphite-checks/graphite_send.py > /dev/null 2>&1 || true
