# cloud/mdb/salt/salt/components/blackbox-sso/conf/supervisor.conf
[program:blackbox-sso]
command=/opt/yandex/mdb-blackbox-sso/bin/mdb-blackbox-sso /etc/blackbox-sso.yaml
process_name=%(program_name)s
directory=/
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=5
user=root
stdout_logfile=/var/log/supervisor/blackbox-sso.log
redirect_stderr=true
