# cloud/salt-formula/common/network/files/90-vrouter.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="vhost0", RUN+="/lib/systemd/systemd-sysctl --prefix=/net/vrouter"