# cloud/mdb/salt/salt/components/downtimer/conf/downtimer.cron
MAILTO=''
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

*/3 *  *   *   *  mdb-downtimer  yazk-flock -c /etc/yandex/mdb-downtimer/zk-flock.json lock 'timeout 180 /etc/cron.yandex/downtimer.py --quiet' 2>&1 >/dev/null || true
