# cloud/mdb/images/dataproc/distro/packages/src/common/hadoop/hadoop-hdfs@.service
[Unit]
Description=Apache Hadoop HDFS %I
Documentation=man:hadoop(1)
After=network.target

[Service]
User=hdfs
Group=hadoop
Environment="DAEMON=%i"
RuntimeDirectory=hadoop-hdfs
ExecStart=/usr/bin/bash -c ' \
  . /etc/default/hadoop && \
  . /etc/default/hadoop-hdfs-$DAEMON && \
  . /usr/lib/bigtop-utils/bigtop-detect-javahome && \
  /usr/bin/hdfs --config /etc/hadoop/conf --daemon start $DAEMON'
ExecStop=/usr/bin/bash -c ' \
  . /etc/default/hadoop && \
  . /etc/default/hadoop-hdfs-$DAEMON && \
  . /usr/lib/bigtop-utils/bigtop-detect-javahome && \
  /usr/bin/hdfs --config /etc/hadoop/conf --daemon stop $DAEMON'
SuccessExitStatus=0 143
WorkingDirectory=/var/lib/hadoop-hdfs
Restart=on-failure
RestartSec=15s
Type=forking

[Install]
WantedBy=multi-user.target
