# cloud/mdb/images/dataproc/distro/packages/src/common/hadoop/hadoop-mapred@.service
[Unit]
Description=Apache Hadoop MapReduce %I
Documentation=man:hadoop(1)
After=network.target

[Service]
User=mapred
Group=hadoop
Environment="DAEMON=%i"
RuntimeDirectory=hadoop-mapreduce
ExecStart=/usr/bin/bash -c ' \
  . /etc/default/hadoop && \
  . /etc/default/hadoop-mapreduce-$DAEMON && \
  . /usr/lib/bigtop-utils/bigtop-detect-javahome && \
  /usr/bin/mapred --config /etc/hadoop/conf --daemon start $DAEMON'
ExecStop=/usr/bin/bash -c ' \
  . /etc/default/hadoop && \
  . /etc/default/hadoop-mapreduce-$DAEMON && \
  . /usr/lib/bigtop-utils/bigtop-detect-javahome && \
  /usr/bin/mapred --config /etc/hadoop/conf --daemon stop $DAEMON'
SuccessExitStatus=0 143
WorkingDirectory=/var/lib/hadoop-mapreduce
Restart=on-failure
RestartSec=15s
Type=forking

[Install]
WantedBy=multi-user.target
