# admins/salt-media/storage/ape/roots/ape-test-storage/etc/cron.d/elliptics_dc
23 4 * * *     root    flock -n /run/lock/elliptics_dc.lock /bin/bash /usr/local/bin/elliptics_dc_cron.sh >/dev/null 2>&1

