# admins/salt-media/storage/ape/roots/ape-test-storage/etc/cron.d/elliptics_defrag
*/10 * * * *     root    flock -n /run/lock/elliptics_defrag.lock /bin/bash /usr/local/bin/elliptics_defrag.sh >/dev/null 2>&1

