# cloud/mdb/salt/salt/components/mdb-dataproc-manager/conf/supervisor.conf
[program:mdb-dataproc-manager]
command=/opt/yandex/dataproc-manager/bin/dataproc-manager --config-path=/opt/yandex/dataproc-manager/etc
environment=GOMAXPROCS=2
directory=/opt/yandex/dataproc-manager
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=mdb-dataproc-manager
stdout_logfile=/var/log/supervisor/mdb-dataproc-manager.log
redirect_stderr=true
