# admins/cult-admin/packages/kassa/config-kassa-mail-fallback/src/etc/cron.d/mail-fallback
*	*	*	*	*   root	lockf -t 1 /var/lock/mail-fallback.lock /usr/bin/mail-fallback.pl >/dev/null 2>&1
