# cloud/mdb/dbaas_infra_tests/images/fake_resourcemanager/config/resourcemanager-mock.supervisor.conf
[program:resourcemanager-mock]
command=/resourcemanager-mock
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/dev/stderr
redirect_stderr=true
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
