# cloud/mdb/salt/salt/components/katan/conf/mdb-katan.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''
LANG=C.UTF-8
LC_ALL=C.UTF-8

# m h dom mon dow user                command
 */5 *  *   *   *  mdb-katan  yazk-flock -c /etc/yandex/mdb-katan/imp/zk-flock.json lock /etc/cron.yandex/mdb-katan-imp.sh >/dev/null 2>&1
 */5 10-18  *   *   *  mdb-katan  yazk-flock -c /etc/yandex/mdb-katan/scheduler/zk-flock.json lock /etc/cron.yandex/mdb-katan-scheduler.sh > /dev/null 2>&1
