# admins/salt-media/storage/elliptics/roots/elliptics-rtc/etc/sv/spacemimic/finish
#!/bin/bash
echo [`date`] `pwd | cut -d/ -f4` >> /var/log/runit_flap.log
ps -ef | grep -v grep | grep "usr/bin/spacemimic" | awk '{print "kill -9 "$2}' | sh

