# cloud/mdb/salt/salt/components/kvm-host/conf/qemu_cleanup.cron
MAILTO=''
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
30 * * * * root sleep $((RANDOM \% 1200)); /etc/cron.yandex/qemu_cleanup.sh >/dev/null 2>/dev/null
