# cloud/mdb/salt/salt/components/common/conf/database_slice_adjuster.cron
MAILTO=''
* * * * * root sleep $((RANDOM \% 55)); flock /tmp/database-slice-adjuster.lock timeout 60 /etc/cron.yandex/database_slice_adjuster.py >/dev/null 2>/dev/null
