# cloud/mdb/salt/salt/components/common/conf/performance-tuner.conf
description "Update cpu frequency governor to performance"

start on mounted MOUNTPOINT=/sys

pre-start script
    test -x /usr/local/yandex/cpu-performance-scaling.sh || { stop; exit 0; }
end script

exec /usr/local/yandex/cpu-performance-scaling.sh
