# admins/git-repo/packages/config-kp-space/debian/config-kp-space/etc/cron.d/seclog-distrub
7 0 * * * root /usr/bin/seclog-distrub.sh
