# admins/salt-master-configs/media_salt_master/common/roots/templates/salt-master/files/csync2/salt-csync2.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
SHELL=/bin/bash
* * * * * root lockf -kst0 -- /tmp/csync2.lock csync2 -x &>> /var/log/csync2.log||true
