# admins/salt-media/noc/roots/units/gitlab/files/etc/cron.d/docker-clean

13 * * * *      root    docker volume prune -f > /dev/null 2>&1; docker system prune -a -f > /dev/null 2>&1

