# admins/salt-media/storage/ape/roots/ape-test-state/etc/cron.d/darkvoice
*/10 * * * *    root    if [ $(ps aux | grep darkvoice | grep config | grep -v grep | grep -v ubic | awk '{print $6}') -gt 13000000 ];then sleep $(( $RANDOM\%120 )); ubic restart darkvoice;fi > /dev/null 2>&1

