# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_innodb_status_output
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
*/15 * * * * mysql flock -n /tmp/my_innodb_status_output_fix.lock /usr/local/yandex/mysql_innodb_status_output.py -d 1800 {{ 'ON' if salt['pillar.get']('data:mysql:config:innodb_status_output', False) else 'OFF' }} >> /var/log/mysql/cron-innodb_status-output.log 2>&1
