# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_autoflush_hosts
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
* * * * * mysql /usr/local/bin/my-autoflush-hosts -t 30s >> /var/log/mysql/cron-my-autoflush-hosts.log 2>&1
