# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_rotate_binlogs
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
0 * * * * mysql /usr/local/bin/my-safe-flush-binlogs -i 1h >> /var/log/mysql/cron-my-safe-flush-binlogs.log 2>&1
