# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_flush_logs
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
* * * * * root flock -n /tmp/mysql_flush_logs.lock /usr/local/yandex/mysql_flush_logs.sh >> /var/log/mysql/cron-my-flush-logs.log 2>&1
