# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_resetup
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
*/5 * * * * root flock -n /tmp/my_resetup.lock -c '/usr/local/bin/my-resetup --resetup-file=/var/run/mysync/mysync.resetup' >> /var/log/mysql/cron-my-resetup.log 2>&1
