# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_mdb_repl_mon
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''
* * * * * mysql flock -n -o /tmp/mysql_mdb_repl_mon.lock /usr/local/yandex/mysql_mdb_repl_mon.py >> /var/log/mysql/cron-mdb-repl-mon.log 2>&1 || true
