# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_rotate_perf_schema
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
0 3 * * * mysql mysql -e "set session lock_wait_timeout = 2; call sys.ps_truncate_all_tables(false)" >> /var/log/mysql/cron-my-rotate-perf-schema.log 2>&1
