# cloud/mdb/salt/salt/components/mysql/conf/cron.d/mysql_purge_binlogs
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
* * * * * mysql /usr/local/bin/my-purge-binlogs -v -s {{ salt['pillar.get']('data:mysql:config:mdb_preserve_binlog_bytes', 1073741824) }} -t 30s >> /var/log/mysql/cron-my-purge-binlogs.log 2>&1
