# admins/salt-media/noc/roots/units/hbf/files/etc/cron.d/update_geo.conf
3 0 * * *	root	p=$(md5sum /etc/nginx/conf.d/geo.conf); /usr/sbin/build_geoconf.py /etc/nginx/conf.d/geo.conf && ([ \"$p\" = \"$(md5sum /etc/nginx/conf.d/geo.conf)\" ] || systemctl reload nginx.service)
