# cloud/mdb/dbaas_infra_tests/images/fake_conductor/config/fake_conductor.conf
[program:fake_conductor]
command=uwsgi /opt/yandex/uwsgi.ini
environment=FAKE_CONDUCTOR_CONFIG=config.py
directory=/opt/yandex
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
user=root
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
