# admins/git-repo/packages/config-caching-dns/debian/postinst
#!/bin/sh

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        /usr/bin/config-caching-dns-update-configs || true

        # CSADMIN-31393 Fix for upgrade from 1.0-53 version
        if [ -f /etc/resolv.conf ] && [ `stat -c '%a' /etc/resolv.conf` -eq 600 ]; then
            chmod 644 /etc/resolv.conf
        fi

        if [ -f /etc/init.d/bind9 ]; then
            if invoke-rc.d bind9 status >/dev/null; then
                invoke-rc.d bind9 stop
            fi
            update-rc.d bind9 disable >/dev/null
        fi

        if ! invoke-rc.d unbound status >/dev/null; then
            invoke-rc.d unbound start
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
