# admins/salt-media/noc/roots/files/nocdev-rs/etc/cron.d/get_vips
17 */6 * * *	root	/etc/network/ya-slb-tun.sh reload eth0
13 */6 * * *       root    curl -s "https://ro.racktables.yandex-team.ru/export/allobjects.php?text=(\{balancer-common\} or \{balancer-nat\}) and not \{в оффлайне\} and not \{disable-natasha-probes\}" | xargs -n 1 host | awk '{print $NF}' > /etc/network/l3_tun
