# cloud/mdb/gpsync/docker/gpsync/supervisor.conf
[program:odyssey]
command=/usr/bin/odyssey /etc/odyssey/odyssey.conf
process_name=%(program_name)s
autostart=false
autorestart=false
stopsignal=TERM
user=gpadmin
stdout_logfile=/var/log/odyssey.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/odyssey.log
stderr_logfile_maxbytes=0

[program:prestart_all]
user=root
command=bash -c "mkdir -p /var/{run,log}/{gpsync,odyssey} && chown gpadmin:gpadmin /var/{run,log}/{gpsync,odyssey}"
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=0
priority=1
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:gpsync]
environment=MASTER_DATA_DIRECTORY="/var/lib/greenplum/data1/master/gpseg-1",GPHOME="/opt/greenplum-db-6",HOME="/home/gpadmin"
command=/usr/local/bin/gpsync -f yes
process_name=%(program_name)s
autostart=false
autorestart=true
stopsignal=TERM
user=gpadmin
priority=10
stdout_logfile=/var/log/gpsync.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/gpsync.log
stderr_logfile_maxbytes=0

[program:greenplum]
environment=MASTER_DATA_DIRECTORY="/var/lib/greenplum/data1/master/gpseg-1",PYTHONPATH="/opt/greenplum-db-6/lib/python",GPHOME="/opt/greenplum-db-6",HOME="/home/gpadmin"
command=/opt/greenplum-db-6/bin/gpstartdisp -a -y -S
process_name=%(program_name)s
autostart=false
autorestart=false
exitcodes=0
startsecs=0
user=gpadmin
priority=10
stdout_logfile=/var/log/greenplum.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/greenplum.log
stderr_logfile_maxbytes=0
