# admins/git-repo/packages/yandex-corba-export/restarters.cron
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-default >/dev/null 2>&1
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-export >/dev/null 2>&1
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-yaca >/dev/null 2>&1
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-yaru >/dev/null 2>&1
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-slovari >/dev/null 2>&1
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-perevod>/dev/null 2>&1
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-weather >/dev/null 2>&1
* * * * *       root    /usr/local/bin/xscript-restart.sh xscript-maps >/dev/null 2>&1
3 * * * *	root	sleep 30; /etc/init.d/lighttpd restart >/dev/null 2>&1
