# admins/git-repo/packages/yandex-corba-export/traffic.cron
* * * * *	root	curl -sf http://export.`hostname -s`.yandex.ru/maps/traffic.xml > /usr/local/www5/export/maps/static_traff.txt
