# admins/git-repo/packages/backend-kinopoisk-common/src/etc/cron.d/generate_indexes
*/5     *       *       *       *  root sleep $((RANDOM\%300)) && flock -w 0 /tmp/generate_indexes.lock /usr/local/bin/generate_indexes.sh 2>&1 > /dev/null
