# cloud/mdb/salt/salt/components/image-releaser/conf/mdb-image-releaser.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''
LANG=C.UTF-8
LC_ALL=C.UTF-8

# m h dom mon dow user                command
{% for db in salt['pillar.get']('data:images') %}
{% set alias = db.get('alias', db.name) %}
 */30 *  *   *   *  mdb-image-releaser  yazk-flock -c /opt/yandex/mdb-image-releaser/zk-flock-{{ alias }}.json lock /etc/cron.yandex/mdb-image-releaser-{{ alias }}.sh
{% endfor %}
