# admins/salt-media/storage/ape/roots/elastic/cron.d/elasticsearch_index_rotate
3 *  * * * elasticsearch INDEX_PREFIX=cocaine-v0.12- KEEP_DAYS=1 KEEP_HOURS=12 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log
5 *  * * * elasticsearch INDEX_PREFIX=noexport- KEEP_DAYS=6 KEEP_HOURS=0 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log
7 *  * * * elasticsearch INDEX_PREFIX=music-test-back- KEEP_DAYS=0 KEEP_HOURS=6 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log
9 *  * * * elasticsearch INDEX_PREFIX=media-test-combainer- KEEP_DAYS=2 KEEP_HOURS=0 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log
11 *  * * * elasticsearch INDEX_PREFIX=loadape-rsyslog- KEEP_DAYS=4 KEEP_HOURS=0 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log
13 *  * * * elasticsearch INDEX_PREFIX=contest- KEEP_DAYS=15 KEEP_HOURS=0 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log
15 *  * * * elasticsearch INDEX_PREFIX=mdst-srw- KEEP_DAYS=4 KEEP_HOURS=0 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log
17 *  * * * elasticsearch INDEX_PREFIX=nginx- KEEP_DAYS=1 KEEP_HOURS=12 timeout 3000 /usr/share/elasticsearch/bin/index_rotate.sh &>>/var/log/elasticsearch/index_rotate.log


# pre-creating indices to avoid rsyslog wait at every new hour
17 *  * * * elasticsearch timeout 300 /usr/share/elasticsearch/bin/index_pre-create.sh &>>/var/log/elasticsearch/index_rotate.log
