# cloud/mdb/dbaas_infra_tests/images/idm-service/docker/idm_service.conf
[program:idm_service]
command=/opt/yandex/mdb-idm-service/mdb-idm-service.wsgi /opt/yandex/mdb-idm-service/uwsgi.ini
environment=IDM_SERVICE_CONFIG=/config/config.py
directory=/tmp
process_name=%(program_name)s
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
