# cloud/mdb/salt/salt/components/deploy/salt-sync/conf/mdb-salt-sync.cron
MAILTO=''
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
LANG=C.UTF-8
LC_ALL=C.UTF-8
SHELL=/bin/sh

*/1 * * * * root zk-flock -c /etc/yandex/mdb-salt-sync/zk-flock.json lock /etc/cron.yandex/mdb-salt-sync.sh 2>&1 >> /var/log/mdb-salt-sync/mdb-salt-sync.log
