# admins/git-repo/packages/config-tools-xfront/cron.d/logster
MAILTO="corba-root@yandex-team.ru"

* * * * *	root	/usr/sbin/logster --parser-option="/etc/logster/logster-xfront.conf" --output=ganglia RPSTimingsLogster /var/log/nginx/access.log 1>/dev/null 2>&1
