# admins/git-repo/packages/config-corba-force-logrotate/force-logrotate.cron

MAILTO=""

0 0 * * *	root	/usr/bin/force-logrotate.sh
