# admins/cult-admin/packages/kassa/yandex-kassa-balancer-autoclose/src/etc/init/balancer-autoclose.conf
# balancer-autoclose - init script for close prestable from balancer
#

#description     "init script for close prestable from balancer"

start on startup

script
     ( /sbin/iptables-save  | grep -s  '\-A INPUT \-s 10.0.0.0\/8 \-p tcp \-m multiport \-\-dports 80\,443 \-j DROP' > /dev/null) || /sbin/iptables -A INPUT -s 10.0.0.0/8 -p tcp -m multiport --dports 80,443 -j DROP
     ( /sbin/ip6tables-save  | grep -s '\-A INPUT \-s fd00\:\:/8 \-p tcp \-m multiport \-\-dports 80\,443 \-j DROP' > /dev/null) ||  /sbin/ip6tables -A INPUT -s fd00::/8 -p tcp -m multiport --dports 80,443 -j DROP
end script
