# cloud/netinfra/rknfilter/yc_rkn_begger_node/etc/cron.d/yc-rkn-begger-node
*/1 * * * * root /usr/bin/yc-rkn-begger-node

