# admins/salt-media/storage/common/roots/units/cocaine/files/etc/cron.d/cocaine-warmup
#
# up cocaine  app workers
#
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/30 * * * * root /usr/local/bin/cocaine-warmup.py >/dev/null 2>/dev/null

