# admins/salt-media/storage/elliptics/roots/files/elliptics-mondb/etc/cron.d/clickhouse-cleaner
13 14 * * *       root    ls /var/lib/clickhouse/data/tacct_db/tacct_nfsense_orig | grep [0-9] | grep -v tmp | cut -d_ -f1 | sort | uniq | while read name; do if [ $(($name + 1)) -lt $(date '+\%Y\%m') ]; then touch '/var/lib/clickhouse/flags/force_drop_table';chown clickhouse:clickhouse '/var/lib/clickhouse/flags/force_drop_table'; echo "ALTER TABLE tacct_db.tacct_nfsense_orig DROP PARTITION $name" | clickhouse-client; fi; done > /dev/null 2>&1

