# admins/salt-media/storage/elliptics/roots/files/elliptics-rs/etc/cron.d/get_vips
17 */6 * * *	root	IF_YA_SLB_TUN="YES" IF_YA_SLB6_TUN="YES" DEBUG="YES" /usr/lib/yandex-netconfig/ya-slb-tun reload eth0
13 */6 * * *       root    curl -s "https://ro.racktables.yandex-team.ru/export/allobjects.php?text=\{балансировщик\} and not \{в оффлайне\} and not \{гипервизор\} and not \{UA балансер\} and not \{балансировщик CDN\} and not \{балансировщик в проектной сети\} and not \{метабалансер\} and not \{unstable\} and not \{\$typeid_1504\} and not \{ОФД\} and not \{балансировщик yacloud\}" | xargs -n 1 host | awk '{print $NF}' > /etc/network/l3_tun
