# admins/salt-media/storage/common/roots/units/atop/files/etc/cron.d/atop
PATH=/bin:/usr/bin:/sbin:/usr/sbin
LOGPATH=/var/log/atop

# daily restart of atop at midnight
0 0 * * * root if [ -d "/run/systemd/system" ]; then systemctl restart atop && find $LOGPATH -name 'atop_*' -mtime +28 -delete; fi
