# cloud/mdb/salt/salt/components/pg-dbs/dbm/conf/dbm_expire_secrets.cron
MAILTO=''
*/15 * * * * postgres sleep $((RANDOM \% 600)); flock -n -w 10 /tmp/dbm_expire_secrets.lock timeout 3600 /etc/cron.yandex/dbm_expire_secrets.py >/dev/null 2>&1
