# admins/salt-media/storage/elliptics/roots/units/mds-hide/files/etc/cron.d/mds-hide-collect-urls
# m   h dom mon dow    user     command
*/4  *   *   *   * s3proxy     /usr/bin/zk-flock mds-hide-collect-urls "/usr/local/bin/collect-ban-urls --config /etc/mds-hide/config.json"
