# cloud/mdb/salt/salt/components/katan/conf/katan/supervisor.conf
[program:mdb-katan]
command=/opt/yandex/mdb-katan/bin/mdb-katan --config-path=/etc/yandex/mdb-katan/katan/
environment=GOMAXPROCS=2
directory=/opt/yandex/mdb-katan
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=mdb-katan
stdout_logfile=/var/log/supervisor/mdb-katan.log
redirect_stderr=true
startretries=1800
