# cloud/mdb/images/dataproc21/distro/packages/src/common/hive/hive-webhcat-server.service
[Unit]
Description=Apache Hive HCatalog WebHCat
Documentation=man:hive(1)
After=network.target

[Service]
User=hive
Group=hive
RuntimeDirectory=hive-hcatalog
ExecStart=/usr/bin/bash -c ' \
  . /etc/default/hadoop && \
  . /usr/lib/bigtop-utils/bigtop-detect-javahome && \
  . /etc/default/hive-hcatalog-server && \
  . /usr/lib/hive-hcatalog/sbin/webhcat_server.sh start'
ExecStop=/usr/bin/bash -c ' \
  . /etc/default/hadoop && \
  . /usr/lib/bigtop-utils/bigtop-detect-javahome && \
  . /etc/default/hive-hcatalog-server && \
  . /usr/lib/hive-hcatalog/sbin/webhcat_server.sh stop'
WorkingDirectory=/var/lib/hive-hcatalog
SuccessExitStatus=0 143
Restart=always
RestartSec=15s
Type=forking

[Install]
WantedBy=multi-user.target
