# nginx/qloud-bin/debians/debian/qloud-nginx.cron.d
# Regular cron jobs for the qloud-nginx
#

MAILTO=""
*/30 *	* * *	root	/usr/sbin/logrotate /etc/logrotate.d/qloud-nginx

