# cloud/mdb/salt/salt/components/dom0porto/conf/dom0_images_sync.cron
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

0 */2 * * * root sleep $((RANDOM \% 3600)) && flock -o -w 10 /tmp/dom0-images-sync.lock timeout 3600 /etc/cron.yandex/dom0-images-sync.sh >>/var/log/dom0-images-sync.log 2>&1
