# cloud/mdb/salt/salt/components/dom0porto/conf/yasmagent-porto-fix.cron
MAILTO=''
SHELL=/bin/bash

*/5 * * * * root /bin/sleep $((RANDOM \% 180)) && (test -x /usr/sbin/portoctl && /usr/sbin/portoctl list -1rt yasmagent_11003 | /bin/grep -q yasmagent_11003 && /usr/sbin/portoctl destroy yasmagent_11003) || /bin/true
