# admins/config-pkgver-hw/debian/config-pkgver-hw.postinst
[ -x /usr/bin/deb-systemd-helper ] || alias deb-systemd-helper=true

##DEBHELPER##
