# admins/salt-media/storage/avatars/roots/files/avatars-proxy/etc/cron.d/postfix_watchdog
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
*/10 * * * *  root sleep `shuf -i 1-60 -n1` && /etc/init.d/postfix status 2>&1 > /dev/null || /etc/init.d/postfix restart 2>&1 > /dev/null
