# admins/salt-media/storage/common/roots/units/walle_agent/files/walle_cms_agent.cron
PATH=/bin:/sbin:/usr/bin:/usr/sbin
* * * * * root flock -n /var/lock/walle_cms_agent.lock /usr/local/bin/agent.py 2>/dev/null
