# admins/salt-master-configs/media_salt_master/common/roots/templates/salt-master/files/arcadia/salt-arc-update.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
SHELL=/bin/bash
*/5 * * * * root sudo -u {{ salt_user }} ARCANUM_TOKEN=`cat /home/{{ salt_user }}/.arc/token` /usr/local/bin/deploy-configs.py -p {{ remote_path }} -t {{ local_target }} --id_rsa /home/{{ salt_user }}/.ssh/id_rsa
