# admins/salt-media/common/roots/templates/tableau/backups/files/tableau_backup.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/tableau/tableau_server/packages/customer-bin.20211.21.0511.0935

0 1 * * 1-6 robot-media-tableau /usr/bin/tableau_backup.py
0 1 * * 0   robot-media-tableau /usr/bin/tableau_backup.py --full
