# admins/salt-media/storage/ape/roots/zipkin/etc/cron.d/zipkin-restart
13 * * * * root sleep $((RANDOM\%30))m && /sbin/restart zipkin-query-service >/dev/null 2>&1

