# cloud/mdb/salt/salt/components/deploy/mdb-deploy-api/conf/mdb-deploy-cleaner.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''
LANG=C.UTF-8
LC_ALL=C.UTF-8

# m h dom mon dow user                command
 0 10 * * WED  mdb-deploy-cleaner  zk-flock -c /etc/yandex/mdb-deploy-cleaner/zk-flock.json lock /etc/cron.yandex/mdb-deploy-cleaner.sh || rm /var/run/mdb-deploy-cleaner/last-exit-status

