# cloud/mdb/datacloud/docker/deploy-salt/config/master.supervisor.conf
[program:salt-master]
command=/usr/bin/salt-master -l debug
autostart=true
autorestart=true
startretries=999999999999
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
