# cloud/mdb/datacloud/docker/deploy-salt/config/update_salt_image.cron
MAILTO=''
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

*/1 * * * * root flock -n -w10 /tmp/update_salt_image.lock /etc/cron.yandex/update_salt_image.py update 2>&1 >/dev/null || true
