# admins/cult-admin/packages/general/config-nginx-cult-graphite-checks/debian/config-nginx-cult-graphite-checks.cron.d
* * * * * root  lockf -t 1 /var/lock/timetail-nginx-log.lock /usr/sbin/timetail-nginx-log.sh
