# cloud/mdb/salt/salt/components/pg-dbs/dbaas_metadb/conf/backup_retry_delete.cron
MAILTO=''
30 */6 * * * postgres sleep $((RANDOM \% 7200)); timeout 3600 /etc/cron.yandex/backup_retry_delete.py >>/var/log/backup-retry-delete.log 2>&1
