# cloud/netinfra/rknfilter/yc_rkn_config_node/etc/cron.d/yc-rkn-config-node
30 12 * * * root /usr/bin/yc-rkn-config-node --dns_resolve_attempts 2 
*/1 * * * * root /usr/bin/yc-rkn-config-node --config_file /etc/yc/rkn/config_node.YaTelecom.yaml --logfile /var/log/yc_rkn_config_node.YaTelecom.log --pid_file_path yc-rkn-config-node.YaTelecom
