# cloud/filestore/packages/yc-nfs-http-proxy-systemd/init/filestore-http-proxy.service
[Unit]
Description=Network File Store HTTP proxy
ConditionPathExists=/Berkanavt/nfs-server/cfg/filestore-http-proxy.txt

[Service]
Restart=always
RestartSec=1
User=nfs
ExecStart=/usr/bin/filestore-http-proxy --config=/Berkanavt/nfs-server/cfg/filestore-http-proxy.txt

[Install]
WantedBy=multi-user.target
