# cloud/mdb/dbaas_infra_tests/images/internal-api/docker_conf/internal_api.conf
[program:internal_api]
command=/opt/yandex/internal-api.wsgi /opt/yandex/uwsgi.ini
environment=DBAAS_INTERNAL_API_CONFIG=/config/config.py
directory=/opt/yandex
process_name=%(program_name)s
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
