# admins/salt-media/storage/elliptics/roots/files/ttl-cleanup/etc/cron.d/ttl-cleanup
34 14 * * * root /usr/bin/ttl-cleanup.py >/dev/null 2>/dev/null
