# admins/salt-media/kp/roots/db/lag/files/check-lag.cron
4 * * * *   root  mysql -e 'SHOW SLAVE STATUS\G' | perl -ne 'if (/SQL_Delay: (\d+)/) { $d = $1; } } { if (defined($d) and ($d != 43200)) { system("mysql", "-e", "STOP SLAVE; CHANGE MASTER TO MASTER_DELAY = 43200; START SLAVE;"); }'
