# cloud/mdb/salt/salt/components/pushclient2/conf/wd-pushclient.cron
MAILTO=''
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

* * * * * root sleep $((RANDOM \% 45)); flock -n /tmp/wd-pushclient.lock timeout 30 /etc/cron.yandex/wd-pushclient.py >/dev/null 2>&1
{% if salt['pillar.get']('data:pushclient:limit_core', 0) != 0 %}
0 * * * * root test -e /usr/local/yandex/log_curator.py && sleep $((RANDOM \% 60)); /usr/local/yandex/log_curator.py clear-old-logs --pool /var/cores --mask 'push-client.*?core$' --minfiles 1 --size 1G
{% endif %}
