# admins/git-repo/packages/corba-yarobot/debian/postinst
#!/bin/bash

update-rc.d yarobot defaults

rm /var/run/yarobot/*

/etc/init.d/yarobot start
