# cloud/bootstrap/api/packages/yc-bootstrap-api/data/lib/systemd/system/yc-bootstrap-api.service
[Unit]
Description=Yandex Cloud Bootstrap API
After=network.target

[Service]
User=yc-bootstrap
ExecStart=/usr/bin/yc-bootstrap-api --config /etc/yc/bootstrap/api/config.yaml

# restart policy
Restart=on-failure
RestartSec=10s
StartLimitBurst=10
StartLimitInterval=1h

[Install]
WantedBy=multi-user.target
