# cloud/mdb/mdb-porto-agent/etc/cron.d/mdb-porto-agent
@daily root /usr/sbin/mdb-porto-agent clean
