# admins/salt-media/storage/elliptics/roots/files/elliptics-l3mon/etc/cron.d/solomonsyncer
13 */6 * * *	root	zk-flock natasha_solomonsyncer -n 1 -w 60 'solomonsyncer -c /etc/solomonsyncer/config.yaml' >/dev/null 2>/dev/null
