# cloud/mdb/salt/salt/components/idm_service/conf/idm_service.conf
[program:mdb_idm_service]
command=/opt/yandex/mdb-idm-service/bin/idm_service.wsgi /opt/yandex/mdb-idm-service/uwsgi.ini
environment=IDM_SERVICE_CONFIG=/opt/yandex/mdb-idm-service/config.py
directory=/opt/yandex/mdb-idm-service
process_name=%(program_name)s
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=web-api
stdout_logfile=/var/log/supervisor/mdb-idm-service.log
redirect_stderr=true
