# cloud/mdb/salt/salt/components/idm_service/conf/rotate_idm_passwords.cron
MAILTO=''
*/5 * * * * web-api sleep $((RANDOM \% 120)); flock -n -w 10 /tmp/idm_pass_rotate.lock curl localhost/rotate-passwords/ 2>&1
