# admins/yandex-media-common-cassandra-backup/src/etc/cron.d/cassandra-backup
# snapshot (full backup)
27 3 * * * root sleep $(($RANDOM / 900)) ; /usr/bin/cassandra-backup-snapshot.sh >> /var/log/cassandra-backup.log 2>&1
