# cloud/mdb/dbaas_infra_tests/images/secrets/config/mdb-secrets.conf
[program:mdb-secrets]
command=/usr/local/sbin/mdb-secrets --config-path /etc --scheme=http --port=8080 --host 0.0.0.0
process_name=%(program_name)s
directory=/
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=5
user=root
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
