# cloud/mdb/salt/salt/components/postgres/conf/index_repack.cron
MAILTO=''
# Index repack
00 23 * * 5 postgres sleep ${RANDOM:0:2}m; flock -n /tmp/pg_index_repack.lock -c "PGAPPNAME=MDB /etc/cron.yandex/repack.py --indices" >/var/log/pg_index_repack.log 2>&1
{% if salt['pillar.get']('data:do_clear_after_repack', True) %}
00 23 * * 4 postgres flock -n /tmp/pg_clear_after_repack.lock -c "/etc/cron.yandex/clear_after_repack.py" >/var/log/pg_clear_after_repack.log 2>&1
{% endif %}
