# cloud/mdb/salt/salt/components/postgres/conf/pgbouncer.conf
{% set bouncer_user = 'postgres' %}
[program:prestart_pgbouncer]
command=/usr/local/bin/prestart_pgbouncer.sh
user=root
group=root
autorestart=false
startsecs=0
priority=1
stdout_logfile=/var/log/supervisor/prestart_pgbouncer.log
stderr_logfile=/var/log/supervisor/prestart_pgbouncer.log

[program:pgbouncer_internal]
process_name=%(program_name)s%(process_num)02d
command=pgbouncer /etc/pgbouncer/pgbouncer_internal%(process_num)02d.ini
numprocs={{ salt['pillar.get']('data:pgbouncer:internal_count', 1) }}
{% if salt['pillar.get']('data:use_pgsync', True) %}
autostart=false
{% endif %}
autorestart=true
stopsignal=TERM
startsecs=0
directory=/etc/pgbouncer/
user={{bouncer_user}}
group={{bouncer_user}}
stdout_logfile=NONE
stderr_logfile=NONE
priority=100

[program:pgbouncer_external]
process_name=%(program_name)s%(process_num)02d
command=pgbouncer /etc/pgbouncer/pgbouncer%(process_num)02d.ini
numprocs={{ salt['pillar.get']('data:pgbouncer:count', 1) }}
{% if salt['pillar.get']('data:use_pgsync', True) %}
autostart=false
{% endif %}
autorestart=true
stopsignal=TERM
startsecs=0
directory=/etc/pgbouncer/
user={{bouncer_user}}
group={{bouncer_user}}
stdout_logfile=NONE
stderr_logfile=NONE
priority=500

[group:pgbouncer_external]
programs=pgbouncer_external

[group:pgbouncer_internal]
programs=pgbouncer_internal
