# cloud/mdb/salt/salt/components/postgres/conf/table_repack.cron
MAILTO=''
# Table repack
0 23 * * 2 postgres sleep ${RANDOM:0:2}m; flock -n /tmp/pg_table_repack.lock -c "PGAPPNAME=MDB /etc/cron.yandex/repack.py --tables --config /etc/pg_table_repack.conf" >/var/log/pg_table_repack.log 2>&1
