# admins/git-repo/packages/config-ganglia-corba-export/cron.d/logster
* * * * *	root	/usr/sbin/logster --parser-option="/etc/logster/lighttpd.conf" --output=ganglia RPSTimingsLogster /var/log/lighttpd-x5/access.log 1>/dev/null 2>&1
* * * * *	root	/usr/sbin/logster --parser-option="/etc/logster/nginx.conf" --output=ganglia RPSTimingsLogster /var/log/nginx/access.log 1>/dev/null 2>&1
