# cloud/mdb/salt/salt/components/pg-dbs/dbaas_metadb/dbaas_cleaner/conf/dbaas-cleaner.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

0 0 * * * root sleep $((RANDOM \% 120)); flock -n -o /tmp/dbaas-cleaner.lock timeout 610m /opt/yandex/dbaas-cleaner/bin/dbaas-cleaner -c /opt/yandex/dbaas-cleaner/config.yaml clean
