# cloud/mdb/deploy/integration_test/images/mdb-deploy-api/config/mdb-deploy-api.supervisor.conf
[program:mdb-deploy-api]
command=/mdb-deploy-api --config-path / --scheme=http --port=8900 --host=::1
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/mdb-deploy-api.log
redirect_stderr=true
