# admins/git-repo/packages/config-tools-balancers-haproxy/debian/postinst
#!/bin/bash

update-rc.d -f haproxy remove
update-rc.d -f haproxy-tools defaults
regenerate-monrun-tasks
