# admins/git-repo/packages/config-tools-haproxy-mysql/debian/config-tools-haproxy-mysql-testing.postinst
#!/bin/bash

update-rc.d -f haproxy remove
/etc/init.d/haproxy stop

if [ -e /etc/init.d/haproxy-tools-mysql ]; then
    update-rc.d -f haproxy-tools-mysql defaults
elif [ -e /etc/init.d/haproxy-tools-mysql-testing ]; then
    update-rc.d -f haproxy-tools-mysql-testing defaults
fi
regenerate-monrun-tasks
