# cloud/mdb/salt/salt/components/datacloud/common/conf/etc/cron.d/apt-cache-cleanup
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

00 04 * * 5 root sleep $((RANDOM \% 120)); apt-get clean >/dev/null 2>&1
