# cloud/mdb/salt/salt/components/datacloud/common/conf/etc/cron.d/cores-cleanup
MAILTO=''
SHELL=/bin/bash

0 4 * * * root sleep $((RANDOM \% 3600)); find /var/cores -ctime +7 -type f -delete >/dev/null 2>&1
