# admins/salt-media/storage/ape/roots/qloud-cloud/etc/sv/graphite-client/run
#!/bin/bash
exec 2>&1
P_H=$(/usr/local/bin/get_qloud_data.py dom0_hostname)
export PORTO_HOST $P_H
touch /var/tmp/graphite_client_ok
exec /usr/local/yabs-graphite-client/graphite-sender.py

