# cloud/mdb/dbaas_worker/dbaas-worker.conf
[program:dbaas-worker]
command=/opt/yandex/dbaas-worker/bin/dbaas-worker -c /etc/dbaas-worker.conf
process_name=%(program_name)s
directory=/opt/yandex/dbaas-worker
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=300
user=root
stdout_logfile=/var/log/supervisor/dbaas-worker.log
redirect_stderr=true
