# cloud/mdb/salt/salt/components/datacloud/dist-sync/conf/mdb-dist-sync.cron
MAILTO=''
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

*/3 *  *   *   *  dist-sync flock --nonblock /tmp/dist-sync.lock /etc/cron.yandex/mdb-dist-sync.py --quiet 2>&1 >/dev/null
