# cloud/mdb/salt/salt/components/docker-host/conf/dbaas_cleanup.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

13 * * * *  root    flock -o -n /tmp/.dbaas_cleanup.lock timeout 1200 /etc/cron.yandex/dbaas_cleanup_hourly.sh >/dev/null 2>&1
43 5 * * *  root    flock -o -n /tmp/.dbaas_cleanup.lock timeout 3600 /etc/cron.yandex/dbaas_cleanup_daily.sh >/dev/null 2>&1
