# admins/salt-media/storage/elliptics/roots/units/yandex-storage-shelf-hw-watcher-conf/src/etc/profile.d/hw-watcher-storage.sh
function watcher-force-disable {
    touch /var/tmp/watcher-force-disable
}
