# cloud/mdb/dbaas_infra_tests/images/go-internal-api/config/go-internal-api.supervisor.conf
[program:mdb-internal-api]
command=/mdb-internal-api --config-path /
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/go-internal_api.log
redirect_stderr=true
