# admins/salt-media/storage/ape_prod/roots/cloud-12/etc/cron.d/actualize_if
*/10 * * * * root /usr/local/bin/actualize_interface.sh

