# admins/salt-media/storage/ape_prod/roots/cloud-12/etc/cron.d/cocaine-runtime-watchdog
# vim: ft=crontab

*/5 * * * *   root   /bin/bash /usr/local/bin/cocaine-runtime-watchdog.sh

