# cloud/mdb/salt/salt/components/walg-server/conf/walg.cron
MAILTO=''
{% for cluster in salt['pillar.get']('data:walg-server:clusters', []) %}
{% set recover_dir = salt['pillar.get']('data:walg-server:recover_dir', '/backups/recover') %}
0 12 * * *	walg	sleep $((RANDOM \% 7200)); flock -n /tmp/.walg_check_{{ cluster }}.lock -c "python3 /etc/cron.yandex/walg_check_backup.py --cluster {{ cluster }} --recover-dir {{ recover_dir }}" >> /var/log/walg/{{ cluster }}.log 2>&1
{% endfor %}
