# cloud/iam/token_agent/package/debian/postinst
#!/bin/bash

case "$1" in
  configure)

    systemctl enable yc-token-agent || :
    systemctl start yc-token-agent || :
  ;;

esac
