# admins/salt-media/noc/roots/units/nocdev-4k/files/etc/cron.d/checkist
# Выгружаем статискику дифов в YT раз в час
0 * * * * root /usr/bin/checkist-yt-export diff_stats --export
0 * * * * root /usr/bin/checkist-yt-export diff_stats --export --acl-safe
