# admins/salt-media/storage/ape/roots/ape-test-cloud-12/etc/cron.d/cocaine-sync-indexes
5,35 * * * *   root  python /usr/local/bin/cocaine-sync-indexes.py >/dev/null 2>&1

