# admins/salt-media/storage/ape/roots/ape-test-cloud-12/etc/cron.d/cocaine-refresh-group
*/5 * * * *   root  /bin/bash /usr/local/bin/cocaine-refresh-group.sh >/dev/null 2>&1

