# admins/salt-media/storage/elliptics/roots/units/wd-mastermind-stat/files/wd-mastermind-stat.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/30 * * * * root /usr/bin/wd-mastermind-stat.sh
