# admins/syslog-ng-watchdog/src/etc/cron.d/syslog-ng-watchdog
#
# restart syslog-ng if it not in process list
#

*/2 * * * *  root  flock -n /tmp/syslog-ng-watchdog.lock /usr/bin/syslog-ng-watchdog.sh &> /dev/null

