# admins/mysql-configurator/debian/mysql-purge-binlogs-4.cron.d
7   *   * * *   root    /usr/bin/mysql-purge-binlogs-4        &>/dev/null || /bin/true

