# cloud/mdb/salt/salt/components/sentry/conf/sentry-cleanup.cron
MAILTO=''
SENTRY_CONF=/etc/sentry
0 3 * * * web-api /opt/yandex/dbaas-sentry/bin/sentry cleanup --days=30 >>/var/log/sentry-cleanup.log 2>&1
