# cloud/mdb/salt/salt/components/mdb-event-producer/conf/supervisor.conf
[program:mdb-event-producer]
command=/opt/yandex/mdb-event-producer/bin/mdb-event-producer
environment=GOMAXPROCS=2
directory=/opt/yandex/mdb-event-producer
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=mdb-event-producer
stdout_logfile=/var/log/supervisor/mdb-event-producer.log
redirect_stderr=true
startretries=1800

