# admins/icecream/backend/debian/postinst
#!/bin/bash

systemctl enable icecream-backend --no-reload || true
