# cloud/salt-formula/common/hw-watcher/files/hw_watcher_walle.cron
HW_WATCHER=/usr/sbin/hw_watcher

10,20,30,40,50 * * * * hw-watcher [ -x $HW_WATCHER ] && $HW_WATCHER disk cron >/dev/null 2>&1

10,20,30,40,50 * * * * hw-watcher [ -x $HW_WATCHER ] && $HW_WATCHER link cron >/dev/null 2>&1

10,20,30,40,50 * * * * hw-watcher [ -x $HW_WATCHER ] && $HW_WATCHER mem cron >/dev/null 2>&1

10,20,30,40,50 * * * * hw-watcher [ -x $HW_WATCHER ] && $HW_WATCHER ecc cron >/dev/null 2>&1
