# cloud/mdb/dbaas_infra_tests/images/salt-master/config/mdb-deploy-saltkeys.supervisor.conf
[program:mdb-deploy-saltkeys]
command=/mdb-deploy-saltkeys
autostart=true
autorestart=true
startretries=999999999999
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/mdb-deploy-saltkeys.log
redirect_stderr=true
