# cloud/mdb/dbaas_infra_tests/images/salt-master/config/nginx.supervisor.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
startretries=999999999999
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/nginx.log
redirect_stderr=true
