# cloud/mdb/dbaas_infra_tests/images/salt-master/config/salt-api.supervisor.conf
[program:salt-api]
command=/usr/bin/salt-api -l debug
autostart=true
autorestart=true
startretries=999999999999
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/salt-api.log
redirect_stderr=true
