# cloud/mdb/salt/salt/components/elasticsearch/conf/restart-kibana.cron
# restart kibana every 3 days
0 3 */3 * * root service kibana restart
