# admins/salt-media/storage/elliptics/roots/elliptics-rtc/etc/sv/elliptics-stat/run
#!/bin/bash

export LC_ALL=en_US.UTF-8
export LANG=en_US.UTF-8
export LANGUAGE=en_US.UTF-8

exec 2>&1
ulimit -c unlimited
sleep 1

exec /usr/bin/elliptics-stat

