# admins/git-repo/packages/config-corba-rasp-wizard-ganglia/cron.d/logster
* *	* * * 	root	/usr/sbin/logster --parser-option="/etc/logster/RPSTimingsLogster.conf" --output=ganglia RPSTimingsLogster /var/log/lighttpd/access.log >/dev/null 2>&1
