# admins/media-graphite-sender/contrib/etc/init/media-graphite-sender.conf
description "Forward various metrics to graphite servers"

start on filesystem or runlevel [2345]
stop on shutdown

setuid media-graphite-sender
setgid media-graphite-sender

respawn
# respawn limit 10 5  # default

exec /usr/bin/media-graphite-sender
