# cloud/mdb/salt/salt/components/redis/common/conf/redis_userfault.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''
LANG=C.UTF-8
LC_ALL=C.UTF-8

# m h  dom mon dow user                command
  * *  *   *   *   root                python /usr/local/yandex/redis_userfault.py >>/var/log/redis/redis-userfault.log 2>&1
