# cloud/mdb/salt/salt/components/common/conf/etc/cron.d/apt-fix-updates
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

43 * * * * root sleep $((RANDOM \% 2700)); timeout 1800 /etc/cron.yandex/apt_fix_updates.py >/dev/null 2>&1
