# cloud/mdb/salt/salt/components/common/conf/etc/cron.d/cores-cleanup
MAILTO=''
SHELL=/bin/bash

0 4 * * * root sleep $((RANDOM \% 3600)); find /var/cores -ctime +7 -type f -delete >/dev/null 2>&1; cd /var/cores; for filename in ./*; do gdb sh "$(echo ${filename: 2})" --eval-command 't a a bt' --batch 2>&1 | grep -e 'is truncated' -e 'is not a core dump' > /dev/null 2>&1 && rm -f "$(echo ${filename: 2})" || true; done;
