# cloud/mdb/salt/salt/components/common/conf/etc/cron.d/errata-apply
MAILTO=''
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

{% set warn_delay = salt.dbaas.errata_warn_threshold() %}

0 10 * * * root sleep $((RANDOM \% 18000)); su monitor -c 'nice -n 20 /usr/local/yandex/monitoring/errata_updates.py -u {{ salt['pillar.get']('data:errata:url', 'https://errata.s3.yandex.net/errata.json.gz') }} -s trapdoor.yandex.net:514/udp 2>/dev/null -w {{ warn_delay }} -c {{ salt['pillar.get']('data:errata:crit_thresh', 30) }}' | grep -v OK && errata_apply.sh
