# cloud/mdb/salt/salt/components/common/conf/etc/cron.d/salt-test-highstate
MAILTO=''
{% set is_dataplane = salt['pillar.get']('data:dbaas:cluster_id') %}
{% set cmd = 'flock -n /tmp/salt_minion_status.lock /usr/bin/timeout 1200 /usr/bin/salt-call state.highstate test=True -l quiet >/dev/null 2>&1' %}
{% if is_dataplane %}
00 5 * * * root sleep $((RANDOM \% 10800)); {{ cmd }}
{% else %}
30 * * * * root sleep $((RANDOM \% 1200)); {{ cmd }}
{% endif %}
