# cloud/mdb/salt/salt/components/common/conf/etc/cron.d/ugly-utmp-fix
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

*/5 * * * * root sleep $((RANDOM \% 120)); if [ $(stat --format=\%s /run/utmp) -gt 8192 ]; then truncate -s 0 /run/utmp ; kill -HUP 1; fi >/dev/null 2>&1
