# cloud/mdb/salt/salt/components/common/conf/etc/cron.d/var-mail-cleanup
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

0 0 4 * * root sleep $((RANDOM \% 120)); rm -f /var/mail/{root,monitor}
