# cloud/mdb/salt/salt/components/common/conf/etc/cron.d/wd-ntp
MAILTO=''
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

*/5 * * * * root sleep $((RANDOM \% 120)); flock -n /tmp/wd-ntp.lock timeout 120 /etc/cron.yandex/wd-ntp.sh >/dev/null 2>&1
