# admins/salt-media/kp/roots/master/files/etc/cron.d/lsyncd-daemon-check
# MANAGED BY SALT
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
* * * * * root (test -f /usr/local/sbin/is-active-master && /usr/local/sbin/is-active-master) || service stop kinopoisk-lsyncd
