# admins/salt-media/kp/roots/master/files/etc/cron.d/sync_last_hour_images_to_s3
0 */1 * * * root /usr/local/sbin/sync_last_hour_images_to_s3.sh
