# admins/salt-media/alet/roots/smartpassproxy/files/etc/supervisor/conf.d/application.conf
[program:application]
command = /app_start.sh
autorestart = true
stopwaitsecs = 10
stopsignal = TERM
