# admins/yandex-conf-hostname-long/src/fixhosts.cron
#
# fix /etc/hosts /etc/hostname (hostname must be in long form)
#
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/15 * * * * root /etc/init.d/fixhosts
