# admins/git-repo/packages/config-init-disable-ra/debian/config-init-disable-ra.cron.d
20 * * * * root flock -x -w 5 /var/lock/disable-ra.lock /etc/init.d/disable-ra restart >/dev/null 2>&1
