# admins/git-repo/packages/config-init-disable-ra/debian/postinst
#!/bin/bash

update-rc.d -f disable-ra defaults

[ -x /usr/sbin/invoke-rc.d ] &&
invoke-rc.d disable-ra start ||
/etc/init.d/disable-ra start
