# admins/git-repo/packages/config-init-disable-ra/debian/prerm
#!/bin/bash

update-rc.d -f disable-ra remove
rm -f /usr/bin/jautov6here-check.sh

echo "Don't forget to set proper RA config for your system now. I kept last state."
