# admins/salt-media/common/roots/templates/mongodb/files/etc/init/mongodbcfg.conf
# Ubuntu upstart file at /etc/init/mongodb.conf

pre-start script
    mkdir -p /var/lib/mongodbcfg/
end script

start on runlevel [2345]
stop on runlevel [06]
respawn

script
  ENABLE_MONGODB="yes"
  if [ -f /etc/default/mongodb ]; then . /etc/default/mongodb; fi
  #Check foreign mongodb.conf
  CONFIG=/etc/mongodbcfg.conf
  if [ -e /etc/mongodbcfg.conf ]; then
    CONFIG=/etc/mongodbcfg.conf
  else
    echo "I can't run Mongodb Config Server without /etc/mongodbcfg.conf"
    exit 0
  fi

  #Mongodb sometimes losts his pid
  dbpath=`cat $CONFIG | egrep ^dbpath | awk -F '=' '{print $2}'`;
  if [ -e $dbpath/mongod.lock ] ; then
	pid=`cat $dbpath/mongod.lock`; 
	cmd=`cat /proc/$pid/cmdline | tr "\000" "\n"|head -n 1 |cut -d : -f 1`
	if [ "$cmd" == "/usr/bin/mongod" ]; then 
	    echo "Mongodb already start"; 
	    exit 1; 
	else 
	    echo "Hmm..old pid.. i remove it"; 
	    rm  $dbpath/mongod.lock; 
	fi; 
    fi

  if [ "x$ENABLE_MONGODB" = "xyes" ]; then exec  start-stop-daemon --start   --pidfile /var/run/mongodbcfg.pid --chuid mongodb --exec  /usr/bin/numactl -- --interleave=all /usr/bin/mongod $SERVER_PARAMS --config $CONFIG; fi
end script
