# admins/salt-media/common/roots/templates/mongodb/files/etc/init/mongos.conf
pre-start script
    mkdir -p /var/log/mongodb/
    chown -R mongodb:mongodb /var/log/mongodb/
end script

start on runlevel [2345]
stop on runlevel [06]
respawn

script
ENABLE_MONGODB="yes"
MPID="/var/run/mongos.pid"
if [ -f /etc/default/mongodb ]; then . /etc/default/mongodb; fi
# change to 1 or 3 comma separated config servers
if [ -e /etc/mongos.conf ]; then
    . /etc/mongos.conf
 else
    echo "I can't run Mongos without /etc/mongos.conf"
    exit 0
fi
if [ -e $MPID ] ; then pid=`cat $MPID`; if `ps aufx| grep $pid | grep -v grep>/dev/null`; then echo "Mongos already start"; exit 1; else echo "Hmm..old pid.. i remove it"; rm $MPID; fi; fi
if [ "x$ENABLE_MONGODB" = "xyes" ]; then exec start-stop-daemon --start --chuid mongodb --pidfile $MPID --exec /usr/bin/numactl -- --interleave=all /usr/bin/mongos $MONGOS_SERVER_PARAMS $RUNARGS; fi
end script
