# cloud/mdb/salt/salt/components/logrotate/conf/logrotate_by_size.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

# Logrotate every minute to enforce size limit.
* * * * * root flock -o -n /tmp/logrotate.lock /etc/cron.yandex/logrotate.sh
