# cloud/mdb/salt/salt/components/logrotate/conf/trim_cores.cron
MAILTO=''
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

# Trim cores if disk space usage is close to the maximum
*/10 * * * * root /usr/local/yandex/log_curator.py clear-old-logs --pool /var/cores --mask 'core' --max-disk-percentage 10 --usage-threshold 95

