# admins/git-repo/packages/xscript-multiple-yt-conf/debian/prerm
#!/bin/bash

set -e

case "$1" in
    remove|upgrade|deconfigure)
        if [ -x "/etc/init.d/xscript-multiple" ]; then
            if [ -x "`which invoke-rc.d >/dev/null 2>&1`" ]; then
                invoke-rc.d xscript-multiple stop yandex-team
            else
                /etc/init.d/xscript-multiple stop yandex-team
            fi
        fi
    ;;

    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
