# admins/salt-media/common/roots/templates/selfdns/files/etc/cron.d/selfdns-client
*/5 * * * * root sleep $((RANDOM\%20))s && /usr/bin/selfdns-client --config /etc/yandex/selfdns-client/media.conf
