# admins/salt-media/noc/roots/units/mysql/files/etc/cron.d/mysql-backup
13 2 * * *	root	grep -vq production /etc/yandex/environment.type || bash /usr/sbin/mysql-backup.sh > /var/log/mysql/backup.log 2>&1
13 1 13 */3 *   root	grep -vq production /etc/yandex/environment.type || bash /usr/sbin/mysql-backup-kvartal.sh > /var/log/mysql/backup.log 2>&1

