# cloud/salt-formula/compute-node-vpc/data-plane/files/vrouter-autorecovery.service
[Unit]
Description='vrouter autorecovery'

[Service]
Type=oneshot
ExecStart=/usr/bin/yc-autorecovery --debug vrouter
