# cloud/salt-formula/compute-node-vpc/data-plane/files/vrouter-autorecovery.timer
[Unit]
Description=Timer for vrouter-autorecovery

[Timer]
OnBootSec={{ boot_delay }}
OnUnitInactiveSec={{ interval }}

[Install]
WantedBy=timers.target
