# admins/salt-media/kp/roots/db/master/files/pt-pause-ptosc.cron
0 22 * * 1-5   root  /usr/local/bin/percona_pause_ptosc.pl --action=pause
0 9 * * 1-5   root  /usr/local/bin/percona_pause_ptosc.pl --action=resume
