# admins/git-repo/packages/config-caching-dns/config-caching-dns.cron
PATH=/usr/bin:/usr/sbin:/sbin:/bin

*/5 * * * * root flock -x -w 5 /var/lock/config-caching-dns-update-configs.lock /usr/bin/config-caching-dns-update-configs >/dev/null 2>&1
