# admins/salt-media/storage/elliptics/roots/units/disk-logbackup/files/etc/cron.d/downloader-logbackup
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

1 3 * * * root /usr/bin/downloader-logbackup.sh >/dev/null
* * * * * root flock -x -w 5 /var/lock/downloads-counter-inc.lock /usr/bin/download_counter.py > /dev/null 2>&1
*/19 * * * * root /bin/bash -c '/usr/bin/timeout -s SIGKILL 15 /usr/bin/get_ticket_for_mpfs_counter.py 2>/dev/null 1>/dev/null' 1>/dev/null 2>/dev/null

