# admins/salt-media/storage/elliptics/roots/files/elliptics-proxy/etc/cron.d/update_noc_allow_rule
35 05 * * * root /usr/local/bin/get_ruleset_mds11020.py

