# cloud/mdb/salt/salt/components/mdb-health/conf/supervisor.conf
[program:mdb-health]
command=/opt/yandex/mdb-health/bin/mdb-health --scheme=http --host=::1 --port=12345 --config-path=/opt/yandex/mdb-health/etc
environment=GOMAXPROCS={{ salt['pillar.get']('data:mdb-health:gomaxprocs', 4) }}
directory=/opt/yandex/mdb-health
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=mdb-health
stdout_logfile=/var/log/supervisor/mdb-health.log
redirect_stderr=true
