# cloud/mdb/deploy/integration_test/images/salt-minion/config/mdb-ping-salt-master.supervisor.conf
[program:mdb-ping-salt-master]
command=/opt/yandex/mdb-ping-salt-master/mdb-ping-salt-master.py
autostart=true
autorestart=true
startretries=999999999999
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/mdb-ping-salt-master.log
redirect_stderr=true
