# admins/git-repo/packages/config-tools-django-fronts/cron.d/logster
* * * * *	root	/usr/sbin/logster --parser-option="/etc/logster/RPSTimingsLogster.conf" --output=ganglia RPSTimingsLogster /var/log/nginx/access.log 1>/dev/null 2>&1

