# cloud/mdb/dbaas_infra_tests/images/dbaas-worker/config/supervisor.conf
[program:dbaas-worker]
command=/opt/yandex/dbaas-worker/bin/dbaas-worker -c /etc/dbaas-worker.conf
process_name=%(program_name)s
directory=/opt/yandex/dbaas-worker
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=300
user=root
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
redirect_stderr=true
