# admins/salt-media/noc/roots/units/cvs/files/etc/cron.d/backup_cvs
13 3 * * *	root	/usr/sbin/backup_cvs.sh >/dev/null 2>&1
