# admins/salt-media/noc/roots/units/cvs/files/etc/cron.d/rsync_from_master
MAILTO="noc-cronjobs@yandex-team.ru"
PATH="/usr/bin:/bin:/usr/sbin:/sbin"
{% if grains['yandex-environment'] == 'testing' %}
*/13	*	*	*	*	root	flock -n /var/lock/rsync_from_master rsync_from_master.sh
{% else %}
*    *       *       *       *       root    flock -n /var/lock/rsync_from_master rsync_from_master.sh
{% endif %}
