# cloud/mdb/salt/salt/components/dbaas-internal-api/conf/internal_api.conf
[program:internal_api]
command=/opt/yandex/dbaas-internal-api/bin/internal-api.wsgi /opt/yandex/dbaas-internal-api/uwsgi.ini
environment=DBAAS_INTERNAL_API_CONFIG=/opt/yandex/dbaas-internal-api/config.py
directory=/opt/yandex/dbaas-internal-api
process_name=%(program_name)s
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=web-api
stdout_logfile=/var/log/supervisor/dbaas-internal-api.log
redirect_stderr=true
