# cloud/mdb/salt/salt/components/dbaas-internal-api/conf/mdb-search-reindexer.cron
MAILTO=''
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

13  2  *   *   *  web-api  yazk-flock -c /etc/yandex/mdb-search-reindexer/zk-flock.json lock '/opt/yandex/dbaas-internal-api/bin/mdb-search-reindexer --config-path=/etc/yandex/mdb-search-reindexer/' >> /var/log/mdb-search-reindexer/search-reindexer.log 2>&1
