# cloud/mdb/salt/salt/components/lazy-trim/conf/lazy-trim.cron
{% set bandwidth = salt['pillar.get']('data:lazy-trim:bandwidth', 256) %}
{% set path = salt['pillar.get']('data:lazy-trim:path', '/data') %}
{% set timeout = salt['pillar.get']('data:lazy-trim:timeout', 36000) %}
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

0 22 * * * root sleep $((RANDOM \% 900)); /opt/yandex/lazy_trim/bin/lazy_trim --bandwidth {{ bandwidth }} --mountpoint {{ path }} --verbose --timeout {{ timeout }} --logfile /var/log/lazy-trim.log
