# admins/git-repo/packages/omniorb-binder/debian/postinst
#!/bin/bash

/usr/sbin/update-rc.d -n omniorb-bindings defaults
