# admins/salt-media/kp/roots/ott/ftpserver/files/etc/cron.d/vsftpd-watchdog
PATH=/usr/sbin:/usr/bin:/sbin:/bin

*/5 * * * * root vsftpd-watchdog.sh &> /dev/null
