# cloud/compute/snapshot/debian/systemd/yc-snapshot.socket
[Socket]
ListenStream=[::]:7627
ListenStream=[::]:7628
ListenStream=[::1]:7629
KeepAlive=true
ReusePort=true

[Install]
WantedBy=sockets.target
