# cloud/mdb/salt/salt/components/dbaas-worker/conf/dbaas-worker-supervisor.conf
[program:dbaas-worker]
command=/opt/yandex/dbaas-worker/bin/dbaas-worker -c /etc/dbaas-worker.conf
process_name=%(program_name)s
directory=/opt/yandex/dbaas-worker
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs={{ salt['pillar.get']('data:dbaas_worker:supervisor_stopwaitsecs', 300) }}
user=web-api
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/dbaas-worker.log
redirect_stderr=true
# disable supervisor's logrotate, cause use system logrotate
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
