# cloud/mdb/salt/salt/components/dbaas-worker/conf/dbaas-worker-agent-cleanup.cron
MAILTO=''
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
*/10 * * * * root sleep $((RANDOM \% 570)); if [ "$(ps aux | grep dbaas-worker | grep -cv grep)" = "1" ]; then pkill -f -u web-api /usr/bin/ssh-agent >/dev/null 2>/dev/null; fi
