# admins/MirrorYandexRu/mirror-scripts/cron/mirror-scripts.cron.d
MAILTO="ashejn@yandex-team.ru"

# 30f79bec32243c31dd91a05c0ad7b80f1e301aea

30 */6 * * * root /usr/bin/adelielinux-linux.sh > /dev/null 2>&1
10 */2 * * * root /usr/bin/almalinux-linux.sh > /dev/null 2>&1
2 */6 * * * root /usr/bin/apache-linux.sh > /dev/null 2>&1
16 */6 * * * root /usr/bin/armbian-apt-linux.sh > /dev/null 2>&1
16 */6 * * * root /usr/bin/armbian-archive-linux.sh > /dev/null 2>&1
16 */6 * * * root /usr/bin/armbian-beta-linux.sh > /dev/null 2>&1
16 */6 * * * root /usr/bin/armbian-dl-linux.sh > /dev/null 2>&1
2 */6 * * * root /usr/bin/atrpms-linux.sh > /dev/null 2>&1
#2 */4 * * * root /usr/bin/agilia-linux.sh > /dev/null 2>&1
13 */4 * * * root /usr/bin/alpine-linux.sh > /dev/null 2>&1
15 7 * * * root /usr/bin/altlinux-beta-linux.sh > /dev/null 2>&1
15 7 * * * root /usr/bin/altlinux-linux.sh > /dev/null 2>&1
12 */6 * * * root /usr/bin/altlinux-nightly-linux.sh > /dev/null 2>&1
12 */6 * * * root /usr/bin/altlinux-nightly-p7-linux.sh > /dev/null 2>&1
30 */1 * * * root /usr/bin/archassault-linux.sh > /dev/null 2>&1
30 */1 * * * root /usr/bin/archlinux-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/archlinuxarm-linux.sh > /dev/null 2>&1
30 */1 * * * root /usr/bin/archlinux32-linux.sh > /dev/null 2>&1
12 */6 * * * root /usr/bin/archserver-linux.sh > /dev/null 2>&1
13 */4 * * * root /usr/bin/astra-linux.sh > /dev/null 2>&1
03 */6 * * * root /usr/bin/astralinux-1.7-linux.sh > /dev/null 2>&1
30 */2 * * * root /usr/bin/blackarch-linux.sh > /dev/null 2>&1
5 */3 * * * root /usr/bin/cassandra-linux.sh > /dev/null 2>&1
5 */3 * * * root /usr/bin/centos-linux.sh > /dev/null 2>&1
5 */3 * * * root /usr/bin/centos-stream-linux.sh > /dev/null 2>&1
20 */2 * * * root /usr/bin/ceph-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/chakra-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/cherokee-linux.sh > /dev/null 2>&1
#10 */6 * * * root /usr/bin/cinnarch-linux.sh > /dev/null 2>&1
3 */3 * * * root /usr/bin/cpan-linux.sh > /dev/null 2>&1
12 */6 * * * root /usr/bin/cygwin-linux.sh > /dev/null 2>&1
20,50 */1 * * * root /usr/bin/dag-linux.sh > /dev/null 2>&1
*/2 * * * * root /usr/bin/debian-backports-linux.sh > /dev/null 2>&1
35 */6 * * * root /usr/bin/debiancd-linux.sh > /dev/null 2>&1
2 3,9,15,21 * * * root /usr/bin/debian-linux.sh > /dev/null 2>&1
*/2 * * * * root /usr/bin/debian-multimedia-linux.sh > /dev/null 2>&1
12 */1 * * * root /usr/bin/debian-ports-linux.sh > /dev/null 2>&1
*/2 * * * * root /usr/bin/debian-security-linux.sh > /dev/null 2>&1
20 */6 * * * root /usr/bin/deepin-packages-linux.sh > /dev/null 2>&1
20 */6 * * * root /usr/bin/deepin-releases-linux.sh > /dev/null 2>&1
3 */6 * * * root /usr/bin/docker-linux.sh > /dev/null 2>&1
10 */6 * * * root /usr/bin/edumandriva-linux.sh > /dev/null 2>&1
12 */4 * * * root /usr/bin/elastic-linux.sh > /dev/null 2>&1
#30 */6 * * * root /usr/bin/elementaryos-linux.sh > /dev/null 2>&1
20,50 */1 * * * root /usr/bin/elrepo-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/epel-linux.sh > /dev/null 2>&1
*/30 * * * * root /usr/bin/fedora-linux.sh > /dev/null 2>&1
*/30 * * * * root /usr/bin/fedora-secondary-linux.sh > /dev/null 2>&1
1 */4 * * * root /usr/bin/fedora-releases-linux.sh > /dev/null 2>&1
7 16 * * * root /usr/bin/fedora-stage-linux.sh > /dev/null 2>&1
2 */12 * * * root /usr/bin/flightgear-linux.sh > /dev/null 2>&1
15 */4 * * * root /usr/bin/freebsd-linux.sh > /dev/null 2>&1
#15 */4 * * * root /usr/bin/ftp.linux.kiev.ua-linux.sh > /dev/null 2>&1
*/30 * * * * root /usr/bin/gentoo-distfiles-linux.sh > /dev/null 2>&1
*/30 * * * * root /usr/bin/gentoo-portage-linux.sh > /dev/null 2>&1
1 */6 * * * root /usr/bin/gnome-linux.sh > /dev/null 2>&1
1 */1 * * * root /usr/bin/gnu-linux.sh > /dev/null 2>&1
1 */2 * * * root /usr/bin/jenkins-linux.sh > /dev/null 2>&1
1 */6 * * * root /usr/bin/kde-linux.sh > /dev/null 2>&1
2 */2 * * * root /usr/bin/kernel-linux.sh > /dev/null 2>&1
30 */12 * * * root /usr/bin/knoppix-dvd-linux.sh > /dev/null 2>&1
2 */12 * * * root /usr/bin/knoppix-linux.sh > /dev/null 2>&1
2 */6 * * * root /usr/bin/kubernetes-linux.sh > /dev/null 2>&1
4 */6 * * * root /usr/bin/libreoffice-linux.sh > /dev/null 2>&1
24 */3 * * * root /usr/bin/ubuntu-cloud-images.sh > /dev/null 2>&1
20 */3 * * * root /usr/bin/launchpad-linux.sh > /dev/null 2>&1
18 4 */7 * * root /usr/bin/mellanox-linux.sh > /dev/null 2>&1
11 */2 * * * root /usr/bin/librealsense.intel.com-linux.sh > /dev/null 2>&1
15 */4 * * * root /usr/bin/linuxmint-linux.sh > /dev/null 2>&1
15 */1 * * * root /usr/bin/linuxmint-packages-linux.sh > /dev/null 2>&1
03 */2 * * * root /usr/bin/lotos-linux.sh > /dev/null 2>&1
20 */4 * * * root /usr/bin/macports-linux.sh > /dev/null 2>&1
*/30 * * * * root /usr/bin/mageia-linux.sh > /dev/null 2>&1
#13 */4 * * * root /usr/bin/mageia-linux-ru-linux.sh > /dev/null 2>&1
30 */6 * * * root /usr/bin/magos-linux.sh > /dev/null 2>&1
40 */3 * * * root /usr/bin/manjaro-linux.sh > /dev/null 2>&1
40 */2 * * * root /usr/bin/mariadb-linux.sh > /dev/null 2>&1
19 5 * * * root /usr/bin/mongodb-linux.sh > /dev/null 2>&1
19 2 * * * root /usr/bin/nodesource-linux.sh > /dev/null 2>&1
39 2 * * * root /usr/bin/nodejs.sh > /dev/null 2>&1
03 */2 * * * root /usr/bin/mcrouter-linux.sh > /dev/null 2>&1
#12 */6 * * * root /usr/bin/msvsphere-linux.sh > /dev/null 2>&1
3 */6 * * * root /usr/bin/mysql-linux.sh > /dev/null 2>&1
12 */3 * * * root /usr/bin/msys2-linux.sh > /dev/null 2>&1
33 */3 * * * root /usr/bin/mxlinux-isos-linux.sh > /dev/null 2>&1
33 */3 * * * root /usr/bin/mxlinux-packages-linux.sh > /dev/null 2>&1
#12 */12 * * * root /usr/bin/naulinux-linux.sh > /dev/null 2>&1
#14 */12 * * * root /usr/bin/naulinux-slce-linux.sh > /dev/null 2>&1
10 */6 * * * root /usr/bin/openbsd-linux.sh > /dev/null 2>&1
5 */1 * * * root /usr/bin/opencsw-linux.sh > /dev/null 2>&1
*/26 * * * * root /usr/bin/redos-linux.sh > /dev/null 2>&1
11 */2 * * * root /usr/bin/ros2-linux.sh > /dev/null 2>&1
5 */1 * * * root /usr/bin/openmandriva-linux.sh > /dev/null 2>&1
30 */1 * * * root /usr/bin/opensuse-linux.sh > /dev/null 2>&1
30 */1 * * * root /usr/bin/opensuse-packman-linux.sh > /dev/null 2>&1
5 */1 * * * root /usr/bin/openvz-linux.sh > /dev/null 2>&1
5 */1 * * * root /usr/bin/osdyson-linux.sh > /dev/null 2>&1
5 */1 * * * root /usr/bin/osdyson-dl-linux.sh > /dev/null 2>&1
5 */2 * * * root /usr/bin/parrot-linux.sh > /dev/null 2>&1
5 */6 * * * root /usr/bin/parabola-linux.sh > /dev/null 2>&1
3 */3 * * * root /usr/bin/percona-linux.sh > /dev/null 2>&1
15 6 * * * root /usr/bin/pointlinux-packages-linux.sh > /dev/null 2>&1
15 6 * * * root /usr/bin/pointlinux-releases-linux.sh > /dev/null 2>&1
27 */2 * * * root /usr/bin/porteus-linux.sh > /dev/null 2>&1
10 */6 * * * root /usr/bin/postgresql-linux.sh > /dev/null 2>&1
15 */3 * * * root /usr/bin/puias-linux.sh > /dev/null 2>&1
23 */3 * * * root /usr/bin/qt-linux.sh > /dev/null 2>&1
15 */6 * * * root /usr/bin/rabbitmq-linux.sh > /dev/null 2>&1
4 */4 * * * root /usr/bin/raspbian-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/redcorelinux-linux.sh > /dev/null 2>&1
10 */2 * * * root /usr/bin/rockylinux-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/rosa-linux.sh > /dev/null 2>&1
1 */2 * * * root /usr/bin/rpmfusion-linux.sh > /dev/null 2>&1
1 4 * * * root /usr/bin/runtu-linux.sh > /dev/null 2>&1
*/2 * * * * root /usr/bin/russianfedora-linux.sh > /dev/null 2>&1
*/2 * * * * root /usr/bin/russianfedora-releases-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/sabayon-linux.sh > /dev/null 2>&1
10 */4 * * * root /usr/bin/sagemath-linux.sh > /dev/null 2>&1
10 */2 * * * root /usr/bin/saltstack-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/scientific-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/siduction-linux.sh > /dev/null 2>&1
2 */12 * * * root /usr/bin/simgear-linux.sh > /dev/null 2>&1
33 */4 * * * root /usr/bin/slackware-linux.sh > /dev/null 2>&1
33 */4 * * * root /usr/bin/slackwarearm-linux.sh > /dev/null 2>&1
15 */1 * * * root /usr/bin/ubuntu-cdimage-linux.sh > /dev/null 2>&1
15 */4 * * * root /usr/bin/ubuntu-cloud-linux.sh > /dev/null 2>&1
11 * * * * root /usr/bin/ubuntu-linux.sh > /dev/null 2>&1
1 */6 * * * root /usr/bin/ubuntu-ports-linux.sh > /dev/null 2>&1
15 */2 * * * root /usr/bin/ubuntu-releases-linux.sh > /dev/null 2>&1
11 */4 * * * root /usr/bin/vlc-linux.sh > /dev/null 2>&1
30 */4 * * * root /usr/bin/voidlinux-linux.sh > /dev/null 2>&1
11 */1 * * * root /usr/bin/xbmc-linux.sh > /dev/null 2>&1

*/2 * * * * root /usr/bin/puppyrus-linux.sh > /dev/null 2>&1
*/2 * * * * root /usr/bin/calculate-linux.sh > /dev/null 2>&1
*/15 * * * * root /usr/bin/opensuse-repositories-linux.sh > /dev/null 2>&1

# clean
30 4 * * * root find /var/log/mirrors/ -type f -mtime +5 -delete > /dev/null 2>&1

13 */4 * * * root /usr/bin/freport.sh > /dev/null 2>&1

# wiki repolist
15 * * * * root /usr/bin/table.sh -f /usr/share/sync-wiki/table.wk > /dev/null 2>&1

# yandexrepo
* * * * * root /usr/bin/yandexrepo.sh > /dev/null 2>&1

