# admins/salt-media/storage/ape_prod/roots/front-12/etc/cron.d/cocaine-runtime
*/3 * * * *	root	if [ `ubic status cocaine-runtime | grep pid | wc -l` -ne 0 ]; then cocaine-tool i > /dev/null; if [ `echo $?` -ne 0 ]; then sleep 5; cocaine-tool i > /dev/null; if [ `echo $?` -ne 0 ]; then sleep 5; cocaine-tool i > /dev/null; if [ `echo $?` -ne 0 ]; then sleep 5; cocaine-tool i > /dev/null; if [ `echo $?` -ne 0 ]; then echo `date` >> /var/log/cocaine-runtime_restart.log; ubic restart cocaine-runtime >> /var/log/cocaine-runtime_restart.log;fi;fi;fi;fi;fi

