# admins/salt-media/storage/ape_prod/roots/front-12/etc/cron.d/cocaine-refresh-group
*/10 * * * *   root    /bin/bash -c 'sleep $(( $RANDOM \% 300 ))' && /bin/bash /usr/local/bin/cocaine-refresh-group.sh >/dev/null 2>&1

