# admins/salt-media/storage/ape_prod/roots/front-12/etc/cron.d/cocaine-unicorn-autofix
# vim: ft=crontab

*/10 * * * *   root rm -f /var/run/cocaine-unicorn-autofix-marker && touch /var/run/cocaine-unicorn-autofix-marker && zk-flock -n 6 cocs_front_autofix -c /etc/distributed-flock-media.json "bash /usr/local/bin/unicorn-autofix.sh" >/dev/null 2>&1

