# admins/salt-media/storage/ape_prod/roots/front-12/etc/cron.d/ua_vpn_list
* * * * * root curl -q http://bobr.s3.mdst.yandex.net/other/proxy_iplist.txt > /tmp/proxy_iplist.txt 2>/dev/null

