# cloud/mdb/salt/salt/components/pg-dbs/pgproxy/conf/supervisor_ping_pgbouncer.conf
[program:ping_pgbouncer]
command = uwsgi -s [::]:8080 --chdir /usr/local/yandex --module webserver_check_pgbouncer:app --protocol=http --processes 4 --harakiri 1
stdout_logfile=/var/log/supervisor/ping_pgbouncer.log
stderr_logfile=/var/log/supervisor/ping_pgbouncer.log
autostart=true
autorestart=true
stopasgroup=true

