# admins/MirrorYandexRu/mirror-scripts/debian/mirror-scripts.cron.d
MAILTO="ashejn@yandex-team.ru"

*/1 * * * * root /usr/bin/clean-log.sh > /dev/null 2>&1
