# admins/salt-media/noc/roots/files/noc-osm/etc/cron.d/ib-mon
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=chernomorets@yandex-team.ru

* * * * * telegraf /opt/ib-mon/ufm-state-check.sh
*/2 * * * * telegraf /opt/ib-mon/update-cache.sh
