# admins/corba-mongo-backup/src/mongo-backup.cron
T_DIR=/var/lib/mongo-backup
# 3600 seconds timeout for mongo-backup enought for all
# timeout need only for prevent run 2 (or more) backups in queue
# problem with mongos stuck backup more then 24 hours (really seen more then 3 days)
12 3 * * * root /usr/bin/lockf -st 3600 /tmp/mongo-oplog-backup.lock /usr/bin/mongo-backup.sh; echo "$? $(date +\%s)" > "$T_DIR/backup.exit"
*/10 * * * * root /usr/bin/lockf -st 0 /tmp/mongo-oplog-backup.lock /usr/bin/mongo-oplog-backup.sh

