# admins/git-repo/packages/corba-cores-sender/sources/etc/cron.d/corba-cores-sender
MAILTO=nemca@yandex-team.ru

* * * * * root /usr/bin/corba-cores-sender.pl
