# admins/salt-media/storage/elliptics/roots/elliptics-rtc/etc/sv/mds-ops-agent/run
#!/bin/bash

export LC_ALL=en_US.UTF-8
export LANG=en_US.UTF-8
export LANGUAGE=en_US.UTF-8

exec 2>&1


exec /usr/local/sbin/mds-ops-agent --http-port 30100 --http-host 127.0.0.1 --http-threads 3

