# admins/salt-media/storage/ape_prod/roots/storage/etc/cron.d/elliptics_merge
*/10 * * * *     root    flock -n /run/lock/elliptics_merge.lock /bin/bash /usr/local/bin/elliptics_merge_cron.sh >/dev/null 2>&1

