# admins/salt-media/noc/roots/files/nocdev-netegress/etc/rc.local
#!/bin/bash
ip link add dummy0 type dummy
sleep 30
systemctl restart networking
sleep 3
ip6tables -t mangle -I POSTROUTING -d 2a02:6b8:0:3400::b0b0 -o eth0 -j IP6SPOOF
ip6tables -t mangle -A OUTPUT -p tcp --src 2a02:6b8::299 -j DSCP --set-dscp-class cs5
iptables -t mangle -A OUTPUT -p tcp --src 87.250.251.42 -j DSCP --set-dscp-class cs5
sleep 30
systemctl restart networking

