# cloud/mdb/salt/salt/components/jenkins/conf/jenkins-backup.cron
AWS_ACCESS_KEY_ID="{{ salt.pillar.get('data:jenkins_backup:key_id') }}"
AWS_SECRET_ACCESS_KEY="{{ salt.pillar.get('data:jenkins_backup:secret_key') }}"
PASSPHRASE="{{ salt.pillar.get('data:jenkins_backup:passphrase') }}"

0 3 * * * root /usr/bin/duplicity -v 4 --asynchronous-upload incr --full-if-older-than 14D /var/backups/jenkins s3://s3.mds.yandex.net/mdb-jenkins-backup/ >>/var/log/jenkins-backup.log 2>&1 && /usr/bin/duplicity -v 4 remove-older-than 30D s3://s3.mds.yandex.net/mdb-jenkins-backup/ >>/var/log/jenkins-backup.log 2>&1
