# admins/salt-media/noc/roots/units/puncher/files/etc/cron.d/puncher-st-watcher-init-only
*/5 * * * * root /usr/bin/puncher-st-watcher --init-only --config=/etc/puncher/puncher.yml && echo "0;OK" > /tmp/st-watcher-monitor || echo "2;Failed" > /tmp/st-watcher-monitor
