# admins/salt-media/noc/roots/units/puncher/files/etc/cron.d/svn-clean
47 3 * * *	puncher	cd /var/spool/puncher/securityfw && /usr/bin/svn cleanup
