# cloud/mdb/salt/salt/components/gerrit/conf/gerrit-backup.cron
AWS_ACCESS_KEY_ID="{{ salt.pillar.get('data:gerrit_backup:key_id') }}"
AWS_SECRET_ACCESS_KEY="{{ salt.pillar.get('data:gerrit_backup:secret_key') }}"
PASSPHRASE="{{ salt.pillar.get('data:gerrit_backup:passphrase') }}"

0 3 * * * root /usr/bin/duplicity -v 4 --asynchronous-upload --exclude '/var/lib/gerrit/logs/*' incr --full-if-older-than 14D /var/lib/gerrit s3://s3.mds.yandex.net/mdb-gerrit-backup/ >>/var/log/gerrit-backup.log 2>&1 && /usr/bin/duplicity -v 4 remove-older-than 30D s3://s3.mds.yandex.net/mdb-gerrit-backup/ >>/var/log/gerrit-backup.log 2>&1
