# cloud/mdb/salt/salt/components/gerrit/conf/supervisor-gerrit.conf
[program:gerrit]
command=/usr/bin/java -Djava.net.preferIPv6Addresses=true -Djava.net.preferIPv6Stack=true -Djava.net.preferIPv4Stack=false -Xms12g -Xmx12g -jar /var/lib/gerrit/bin/gerrit.war daemon -d /var/lib/gerrit
process_name=%(program_name)s
directory=/var/lib/gerrit
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=300
user=gerrit
stdout_logfile=/var/log/supervisor/gerrit.log
redirect_stderr=true
