# admins/salt-media/storage/elliptics/roots/files/elliptics-stat/etc/cron.d/logstore-rotator
SHELL=/bin/bash
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=mail-root@yandex-team.ru

01 4 * * * root flock -n /tmp/logstore-rotator-sas.lock -c '/usr/sbin/logstore-rotator.py /etc/logstore-rotator.conf sas 2>&1'
01 5 * * * root flock -n /tmp/logstore-rotator-man-iva.lock -c '/usr/sbin/logstore-rotator.py /etc/logstore-rotator.conf man-iva 2>&1'
01 6 * * * root flock -n /tmp/logstore-rotator-myt.lock -c '/usr/sbin/logstore-rotator.py /etc/logstore-rotator.conf myt 2>&1'
01 7 * * * root flock -n /tmp/logstore-rotator-ugr-fol.lock -c '/usr/sbin/logstore-rotator.py /etc/logstore-rotator.conf ugr-fol 2>&1'
