# cloud/mdb/salt/salt/components/dbaas-cron/conf/wd-dbaas-cron.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * root sleep $((RANDOM \% 55)); /etc/cron.yandex/wd-dbaas-cron.sh >/dev/null 2>&1
