# admins/git-repo/packages/yandex-kp-hsdb-clean/src/etc/cron.d/hs-clear
30 4 * * * root /usr/bin/hs-clear.sh >> /var/log/hs-clear.log 2>&1
