# cloud/mdb/salt/salt/components/mdb-salt-sync/conf/mdb-salt-sync.cron
MAILTO=''
SHELL=/bin/sh

*/10 * * * * root cd /opt/yandex/mdb-salt-sync && flock -n /tmp/mdb-salt-sync.lock timeout 600 ./full-sync-and-distribute.sh >>/var/log/mdb-salt-sync.log 2>&1
