# cloud/vpc/images/salt/mr-prober/creator/files/meeseeks-updater.timer
[Unit]
Description=Run Meeseeks Compute Node List Updater Every Minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
