# cloud/mdb/salt/salt/components/redis/server/conf/redis_rewrite_aof.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''
LANG=C.UTF-8
LC_ALL=C.UTF-8

# m h  dom mon dow user                command
  * *  *   *   *   root                /usr/local/yandex/redis_rewrite_aof.py >>/var/log/redis/redis-rewrite-aof.log 2>&1; /usr/local/yandex/redisctl/redisctl.py --action remove_aof --dry-run
