# cloud/mdb/dbaas_infra_tests/images/fake_iam/config/access_service.conf
[program:access_service_mock]
command=python3 -m yc_as_mock
directory=/
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
user=nobody
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
