# cloud/netinfra/rknfilter/yc_rkn_filter_node/etc/cron.d/yc-rkn-filter-node
# for filter nodes in prod/preprod
*/1 * * * * root /usr/bin/yc-rkn-filter-node
# for filter nodes in Yatelecom
#*/1 * * * * root /usr/bin/yc-rkn-filter-node --config_file /etc/yc/rkn/filter_node.YaTelecom.yaml
