# admins/salt-media/storage/ape_prod/roots/ape-state/etc/cron.d/srw_runlist
* * * * *	root	flock -n /run/lock/srw.lock python /usr/local/bin/srw_runlist.py >/dev/null 2>&1

