# cloud/mdb/salt/salt/components/mysql/conf/mysync.upstart.conf
description "MySync HA tool for MySQL"
author "mialinx"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
    mkdir -p /var/run/mysync
    chown mysql:mysql /var/run/mysync
end script

respawn

exec setuidgid mysql /usr/local/bin/mysync --config /etc/mysync.yaml --loglevel=Debug
