# admins/salt-media/kp/roots/db/binary/files/etc/cron.d/mysql-backup
*/5 3-4 * * 3   root  /usr/bin/mysql-backup-data-4    &>/dev/null || true
