# cloud/mdb/dbaas_infra_tests/images/fake_smtp/config/fake_smtp.conf
[program:fake_smtp]
command=python3 /opt/yandex/fake_smtp.py
directory=/opt/yandex
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
user=root
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
