# admins/salt-media/noc/roots/units/nocdev-grads/files/etc/systemd/system/openvpn@.service.d/override.conf
# default LimitNPROC not works in an unprivileged container
[Service]
LimitNPROC=infinity
