# cloud/blockstore/packages/yc-nbs-http-proxy-systemd/init/blockstore-http-proxy.service
[Unit]
Description=Network Block Storage HTTP proxy
ConditionPathExists=/Berkanavt/nbs-server/cfg/nbs-http-proxy.txt

[Service]
Restart=always
RestartSec=1
User=nbs
ExecStart=/usr/bin/blockstore-http-proxy --config=/Berkanavt/nbs-server/cfg/nbs-http-proxy.txt

[Install]
WantedBy=multi-user.target
