# cloud/mdb/salt/salt/components/dbaas-compute/conf/mdbuild.cron
MAILTO=''
@reboot root /usr/local/sbin/mdbuild.sh >> /var/log/mdbuild.log 2>&1
