# cloud/mdb/salt/salt/components/mdb-maintenance/cron/maintenance-sync.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
MAILTO=''

5 * * * * root yazk-flock -c /etc/yandex/mdb-maintenance/zk-flock.json lock /etc/cron.yandex/mdb-maintenance-sync.sh || rm /var/run/mdb-maintenance-sync/last-exit-status
