# cloud/mdb/images/dataproc/distro/packages/src/common/hive/hive-metastore.service
[Unit]
Description=Apache Hive Metastore
Documentation=man:hive(1)
After=network.target

[Service]
User=hive
Group=hive
RuntimeDirectory=hive
StandardOutput=null
StandardError=null
SyslogIdentifier=hive-metastore
ExecStart=/usr/bin/bash -c ' \
  . /etc/default/hadoop && \
  . /etc/default/hive-metastore && \
  . /usr/lib/bigtop-utils/bigtop-detect-javahome && \
  /usr/lib/hive/bin/hive --config /etc/hive/conf --service metastore'
WorkingDirectory=/var/lib/hive
Restart=always
RestartSec=15s
Type=simple

[Install]
WantedBy=multi-user.target
