# admins/git-repo/packages/syslog-ng-watchdog-corba/src/etc/cron.d/syslog-ng-watchdog-corba
#
# restart syslog-ng if it not in process list
#

PATTERN="\/sbin\/syslog-ng"
RESTART_CMD="/etc/init.d/syslog-ng restart"
SLEEP_DELAY="10"

*/2 * * * *  root  if ! ps ax | grep -E "$PATTERN" | grep -v grep >/dev/null; then sleep $SLEEP_DELAY; if ! ps ax | grep -E "$PATTERN" | grep -v grep >/dev/null; then $RESTART_CMD >/dev/null 2>&1; fi ; fi
