# cloud/mdb/salt/salt/components/clickhouse/conf/clickhouse-server.cron
MAILTO=''
* * * * * root /usr/local/yandex/clickhouse-wd.sh 1>/dev/null 2>&1
