# cloud/mdb/salt/salt/components/clickhouse/conf/log-curator.cron
MAILTO=''
# Truncate clickhouse logs if logfile too small for logrotate and prerotate section won't run

0 * * * * root flock -n -o /tmp/log-curator.lock /usr/local/yandex/log_curator.py truncate-unrotatable --pool /var/log/clickhouse-server/ --mask 'clickhouse-server.*log'
