# admins/yandex-media-common-redis-slave-watcher/debian/cron.d
0 * * * * root [ -e /etc/yandex/redis-slave-watcher.yaml ] && /usr/bin/redis-slave-watcher
