# admins/salt-media/storage/elliptics/roots/elliptics-rtc/etc/sv/mecoll/run
#!/bin/bash
exec 2>&1
exec python -m mecoll --app-tag tier

