# cloud/salt-formula/slb-adapter/files/dns-check.cron
*/5 * * * * root /usr/local/bin/dns-check.sh

