# cloud/salt-formula/slb-adapter/files/wd-nginx.cron
*/1 * * * * root pgrep --count --exact --full 'nginx: master process /usr/sbin/nginx -c /etc/nginx/nginx.conf' >/dev/null || systemctl restart nginx

