# admins/salt-media/storage/ape/roots/qloud-cloud/etc/sv/cocaine-refresh-group/run
#!/bin/bash
while [ true ]
do
    /bin/bash /usr/local/bin/cocaine-refresh-group.sh
    sleep 300
done

