# cloud/mdb/dbaas_infra_tests/images/mdb-deploy-api/config/nginx.supervisor.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/nginx.log
redirect_stderr=true
