# cloud/mdb/mdb-config-salt/src/mdb-salt-returner.cron
* * * * * root sleep $(( RANDOM \% 55 )); flock -n /tmp/mdb-salt-returner.lock timeout 180 /usr/bin/python2 /usr/local/yandex/salt/returners/mdb-salt-returner.py >>/var/log/mdb-salt-returner.log 2>&1
