# admins/pingunoque/debian/pingunoque.upstart
# pingunoque daemon

description    "pingunoque daemon"

start on runlevel [345]
stop on runlevel [06]

pre-start script
    /usr/local/bin/pingunoque_prepare.sh start

end script

pre-stop script
    /usr/local/bin/pingunoque_prepare.sh stop
end script

exec /usr/bin/pingunoque
