# cloud/mdb/deploy/integration_test/images/salt-master/config/ipv6_etc_hosts.supervisor.conf
[program:ipv6_etc_hosts]
command=/ipv6_etc_hosts.sh
autostart=true
autorestart=false
startretries=1
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/ipv6_etc_hosts.log
redirect_stderr=true
