# cloud/mdb/deploy/integration_test/images/salt-master/config/master.supervisor.conf
[program:salt-master]
command=/usr/bin/salt-master -l debug
autostart=true
autorestart=true
startretries=999999999999
stopasgroup=true
stopsignal=TERM
user=root
stdout_logfile=/var/log/supervisor/salt-master.log
redirect_stderr=true
