# admins/salt-media/storage/ape/roots/ape-test-front-12/etc/cron.d/ape-echo-reqtime
# vim: ft=crontab

* * * * *   root  /bin/bash /usr/local/bin/check-cocaine-reqtime.sh

