# admins/salt-media/storage/ape/roots/ape-test-front-12/etc/cron.d/cocaine-crashlog-clean
# vim: ft=crontab

3 5 * * *	root	 sleep $(( RANDOM \% 6000));cocaine-tool crashlog cleanrange --from_day 13-07-2010 --up_to_day `date -d -3day +\%d-\%m-\%Y` --timeout 100000 >/dev/null 2>&1

