# admins/salt-media/storage/elliptics/roots/units/mastermind/files/etc/cron.d/config-mds-mastermind
12,42 * * * * root /usr/bin/zk-flock elliptics-collector_mastermind_cache_clean "/usr/bin/mastermind cache clean" >/dev/null 2>&1
