# admins/salt-media/storage/elliptics/roots/units/s3/files/etc/cron.d/s3-dispenser-report-usage
# m   h dom mon dow    user     command
*/10  *   *   *   * s3proxy     /usr/bin/zk-flock s3-abc-dispenser-report-usage "/usr/bin/s3-abc-dispenser-report-usage --config /etc/s3-dispenser/s3-dispenser.json"
5     *   *   *   * s3proxy     /usr/bin/zk-flock s3-abc-dispenser-report-usage "/usr/bin/s3-abc-dispenser-report-usage --config /etc/s3-dispenser/s3-dispenser.json --set-quotas"
